# www-misc/vdradmin-am/files/vdradmind.service
[Unit]
Description=WWW interface daemon for VDR
Documentation=man:vdradmind(8)
Documentation=file:///usr/share/doc/vdradmin-am/INSTALL

[Service]
ExecStart=/usr/bin/vdradmind --nofork
User=vdradmin
StandardOutput=syslog
SyslogLevel=warning

[Install]
WantedBy=multi-user.target
