# sys-apps/usbguard/files/usbguard-0.7.6-usbguard-dbus.openrc
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/usbguard-dbus"
command_args="--system"
command_background="true"
pidfile="/run/usbguard-dbus.pid"

depend() {
    need dbus usbguard
}
