# sys-cluster/keepalived/files/init-keepalived
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

opts="reload"

depend() {
	use logger
	# The interfaces do not actually need to exist to start, it handles them gracefully.
	use net 
}

extra_commands="checkconfig"

PIDFILE=/var/run/keepalived.pid

checkconfig() {
	# keepalived has a config check command, but it does not work while the daemon is running!
	if [ ! -e /etc/keepalived/keepalived.conf ] ; then
		eerror "You need an /etc/keepalived/keepalived.conf file to run keepalived"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting Keepalived"
	start-stop-daemon --start --quiet --pidfile $PIDFILE \
		--exec /usr/sbin/keepalived -- $opts
	eend $?
}

stop() {
	ebegin "Stopping Keepalived"
	start-stop-daemon --stop --quiet --pidfile $PIDFILE
	eend $?
}

reload() {
	ebegin "Reloading keepalived.conf"
	start-stop-daemon --pidfile $PIDFILE --signal HUP
	eend $?
}
