# www-servers/boa/files/boa.service
[Unit]
Description=Boa web server
After=network.target

[Service]
ExecStart=/usr/sbin/boa -d

[Install]
WantedBy=multi-user.target
