# net-misc/babeld/files/babeld
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

description="Babel routing daemon"
command=babeld
pidfile=/run/$command.pid
command_args="-I $pidfile -S /var/lib/babeld/state"
command_args_background="-D"
retry="TERM/30/KILL/5"
