# sys-power/nvclock/files/nvclock_initd
#!/sbin/runscript

depend() {
	need xdm
}

start() {
	ebegin "Setting nvidia card overclock speed"
	nvclock ${CLOCKOPTS} >/dev/null
	eend $?
}

