# sys-power/nvclock/files/nvclock_initd-r1
#!/sbin/openrc-run

depend() {
	need display-manager
}

start() {
	ebegin "Setting nvidia card overclock speed"
	nvclock ${CLOCKOPTS} >/dev/null
	eend $?
}

