# net-nds/yp-tools/files/domainname.service
[Unit]
Description=NIS Domainname

[Service]
Type=oneshot
ExecStart=/usr/bin/nisdomainname $NISDOMAINNAME
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
