# net-p2p/go-ipfs/files/ipfs.service
[Unit]
Description=InterPlanetary File System
After=network.target

[Service]
CPUSchedulingPolicy=batch
Environment=IPFS_PATH=/var/lib/ipfs/.ipfs
ExecStart=/usr/bin/ipfs daemon --enable-gc --migrate
ExecStop=/usr/bin/ipfs shutdown
Group=ipfs
LimitNOFILE=8192
PrivateTmp=yes
ProtectSystem=full
Restart=on-failure
Type=simple
User=ipfs

[Install]
WantedBy=multi-user.target
