# net-misc/gofish/files/gofish-1.2-r1.rc
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/gofish"
command_args="${GOFISH_START_OPTS}"
command_args_background="-d"
pidfile="/var/run/gofish.pid"

depend() {
       need net
}
