# sys-fabric/rdma-core/files/iwpmd.init
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/iwpmd"
pidfile="/run/${SVCNAME}.pid"
start_stop_daemon_args="-b -m"
