# sys-devel/icecream/files/iceccd.service
[Unit]
Description=Icecream Distributed Compiler
After=network.target nss-lookup.target

[Service]
Type=simple
SyslogIdentifier=iceccd
ExecStart=/usr/sbin/iceccd -u icecream -b /var/cache/icecream -l /var/log/icecream/iceccd -v
Nice=5

[Install]
WantedBy=multi-user.target
