# net-misc/tinyssh/files/tinyssh.socket
[Unit]
Description=TinySSH service (socket-activated)
Documentation=https://tinyssh.org
Wants=tinyssh-makekey.service

[Socket]
ListenStream=%i
Accept=true
KeepAlive=true
IPTOS=low-delay

[Install]
WantedBy=sockets.target
