# dev-perl/Mail-SpamAssassin-Plugin-GoogleSafeBrowsing/files/update_google_safebrowsing.cron
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-perl/Mail-SpamAssassin-Plugin-GoogleSafeBrowsing/files/update_google_safebrowsing.cron,v 1.1 2008/07/05 09:46:28 robbat2 Exp $
#
*/30 * * * *	root	/usr/sbin/update_google_safebrowsing.sh
