# net-nntp/leafnode/files/texpire.cron
#!/bin/sh
#
# Uncomment the following if you want to
# clean up your news spool regularly.
# 
#
# This is part of the leafnode package.
#

#/usr/sbin/texpire
