# x11-drivers/ati-drivers/files/atieventsd.service
[Unit]
Description=Catalyst event Daemon
Requires=acpid.socket

[Service]
ExecStart=/opt/sbin/atieventsd --nodaemon

[Install]
WantedBy=multi-user.target
