# app-emulation/simplevirt/files/simplevirtd.initd
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="simplevirtd"
pidfile="/run/simplevirtd.pid"
retry="SIGTERM/7200"
command="/usr/sbin/simplevirtd"
command_background="1"
command_args="--syslog ${SIMPLEVIRTD_ARGS}"

depend() {
	need net
	use dns logger
}
