# dev-lang/erlang/files/epmd.service
[Unit]
Description=Erlang Portmapper Daemon
Wants=network.target

[Service]
ExecStart=/usr/bin/epmd

[Install]
WantedBy=multi-user.target
