# net-print/cups-filters/files/cups-browsed.init.d-r1
#!/sbin/openrc-run

pidfile="/var/run/cups-browsed.pid"
command="/usr/sbin/cups-browsed"
command_background="true"

depend() {
	need cupsd avahi-daemon
}
