# net-proxy/dante/files/dante-sockd.service
[Unit]
Description=SOCKS v4 and v5 compatible proxy server and client
After=network.target

[Service]
ExecStart=/usr/sbin/sockd 

[Install]
WantedBy=multi-user.target
