# mail-filter/spampd/files/init
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

# Note: spampd configuration file is /etc/conf.d/spampd

depend() {
	use net
}

start() {
	ebegin "Starting spampd"
	start-stop-daemon --start --quiet \
		--exec /usr/sbin/spampd -- --pid=/var/run/spampd.pid \
			${SPAMPD_OPTS}
	eend $? "Failed to start spampd"
}

stop() {
	ebegin "Stopping spampd"
	start-stop-daemon --stop --quiet --pidfile /var/run/spampd.pid
	eend $? "Failed to stop spampd"
}
