# app-i18n/freewnn/files/freewnn.initd
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/Wnn4/jserver"

depend() {
	use logger
	need net
}

stop() {
	ebegin "Stopping ${name:-${RC_SVCNAME}}"
	/usr/bin/Wnn4/wnnkill
	eend ${?}
}
