# app-i18n/freewnn/files/freewnn.initd-r1
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/jserver"
command_args="${FREEWNN_OPTS}"

depend() {
	use logger
	need net
}

stop() {
	ebegin "Stopping ${name:-${RC_SVCNAME}}"
	/usr/sbin/wnnkill
	eend ${?}
}
