# app-misc/tmux/files/tmux.socket
[Socket]
ListenStream=/tmp/tmux-%U/%i
DirectoryMode=0700
SocketMode=0660

[Install]
WantedBy=sockets.target
DefaultInstance=default
