# app-misc/tmux/files/tmux.service
[Unit]
Requires=tmux@%i.socket

[Service]
ExecStart=/usr/bin/tmux -D

[Install]
Also=tmux@%i.socket
DefaultInstance=default
