# app-admin/glance/files/glance.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts ${SVCNAME} service for OpenStack"

command=/usr/bin/${RC_SVCNAME}
command_user="${GLANCE_USER:-glance}"
command_background=yes
pidfile=/run/${RC_SVCNAME}.pid
required_files=/etc/glance/${RC_SVCNAME}.conf

depend() {
	need net
}
