# app-admin/glance/files/glance.initd-r2
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts ${RC_SVCNAME} service for OpenStack"

command=/usr/bin/${RC_SVCNAME}
command_user="${GLANCE_USER:-glance}"
command_background=yes
pidfile=/run/${RC_SVCNAME}.pid
required_files=/etc/glance/${RC_SVCNAME}.conf

start_stop_daemon_args="--quiet"

depend() {
	need net
}
