# net-misc/cfengine/files/cf-execd.rc6
#!/sbin/openrc-run

depend() {
	need net
	after cf-serverd
}

start() {
	ebegin "Starting cf-execd"
	start-stop-daemon --start --exec /usr/sbin/cf-execd -- ${EXTRA_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping cf-execd"
	start-stop-daemon --stop --pidfile /var/cfengine/cf-execd.pid --exec /usr/sbin/cf-execd
	eend $?
}
