# net-misc/cfengine/files/cfenvd.rc6
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting cfenvd"
	start-stop-daemon --start --quiet --exec /usr/sbin/cfenvd
	eend $?
}

stop() {
	ebegin "Stopping cfenvd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/cfenvd
	eend $?
}
