# sys-fs/btrfsmaintenance/files/btrfs-scrub.timer
[Unit]
Description=Scrub btrfs filesystem, verify block checksums
Documentation=man:btrfs-scrub

[Timer]
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
