# sys-fs/btrfsmaintenance/files/btrfsmaintenance-refresh.service
[Unit]
Description=Update cron periods from /etc/sysconfig/btrfsmaintenance
After=local-fs.target

[Service]
ExecStart=/usr/share/btrfsmaintenance/btrfsmaintenance-refresh-timers.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
