# sys-fs/btrfsmaintenance/files/btrfs-scrub.service
[Unit]
Description=Scrub btrfs filesystem, verify block checksums
Documentation=man:fstrim
After=fstrim.service btrfs-trim.service

[Service]
Type=oneshot
ExecStart=/usr/share/btrfsmaintenance/btrfs-scrub.sh
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
