# sys-fs/btrfsmaintenance/files/btrfs-defrag.service
[Unit]
Description=Defragment file data on a mounted filesystem
Documentation=man:btrfs-filesystem
After=fstrim.service btrfs-trim.service btrfs-scrub.service

[Service]
Type=oneshot
ExecStart=/usr/share/btrfsmaintenance/btrfs-defrag.sh
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
