# app-admin/monit/files/monit.initd-5.0
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/app-admin/monit/files/monit.initd-5.0,v 1.2 2011/12/31 16:28:50 idl0r Exp $

extra_started_commands="reload"

depend() {
	use net
}

start() {
	ebegin "Starting monit"
	start-stop-daemon --start --quiet --exec /usr/bin/monit >/dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping monit"
	start-stop-daemon --stop --quiet --exec /usr/bin/monit
	eend $?
}

reload() {
	ebegin "Reloading monit"
	/usr/bin/monit -c /etc/monitrc reload >/dev/null 2>&1
	eend $?
}
