# net-analyzer/greenbone-security-assistant/files/gsad
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/greenbone-security-assistant/files/gsad,v 1.2 2012/03/27 11:51:13 scarabeus Exp $ 

depend() {
	need net
}

start() {
	ebegin "Starting Greenbone Security Assistant Daemon"
	start-stop-daemon --start --name gsad \
		--exec /usr/sbin/gsad \
		--pidfile /var/run/gsad.pid
	eend $?
}

stop() {
	ebegin "Stop openvassd"
	start-stop-daemon --stop --name gsad \
		--pidfile /var/run/gsad.pid
	eend $?
}
