# net-analyzer/greenbone-security-assistant/files/gsad.service
[Unit]
Description=OpenVAS Manager
After=network.target
After=openvassd.service
After=openvasmd.service
Wants=openvasmd.service

[Service]
Type=forking
EnvironmentFile=-/etc/openvas/sysconfig/gsad-daemon.conf
ExecStart=/usr/sbin/gsad $OPENVAS_SECURITY_ASSISTANT_OPTIONS $OPENVAS_SECURITY_ASSISTANT_LISTEN_ADDRESS $OPENVAS_SECURITY_ASSISTANT_LISTEN_PORT $OPENVAS_SECURITY_ASSISTANT_MANAGER_LISTEN_ADDRESS $OPENVAS_SECURITY_ASSISTANT_MANAGER_LISTEN_PORT $OPENVAS_SECURITY_ASSISTANT_GNUTLS_PRIORITIES
User=root
Group=root
TimeoutSec=1200

[Install]
WantedBy=multi-user.target
