# net-analyzer/greenbone-security-assistant/files/gsad.service
[Unit]
Description=Greenbone Security Assistant Daemon
After=network.target
Before=openvasmd.service

[Service]
EnvironmentFile=-/etc/openvas/gsad-daemon.conf
ExecStart=/usr/sbin/gsad -f ${ASSISTANT_LISTEN} ${ASSISTANT_PORT} ${MANAGER_LISTEN} ${MANAGER_PORT} ${ASSISTANT_EXTRA_ARGS}
Restart=always
RestartSec=1
User=root
Group=root
TimeoutSec=1200

[Install]
WantedBy=multi-user.target
