# net-analyzer/greenbone-security-assistant/files/gsad
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting Greenbone Security Assistant Daemon"
	start-stop-daemon --start --name gsad \
		--exec /usr/sbin/gsad \
		--pidfile /var/run/gsad.pid
	eend $?
}

stop() {
	ebegin "Stop openvassd"
	start-stop-daemon --stop --name gsad \
		--pidfile /var/run/gsad.pid
	eend $?
}
