# net-misc/anydesk/files/anydesk.init
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="AnyDesk service"
description="AnyDesk global service"
command="/opt/bin/anydesk"
command_args="--service"
pidfile="/run/anydesk.pid"
command_background=true
