# net-misc/anydesk/files/anydesk-4.0.1.service
[Unit]
Description=AnyDesk
Requires=network.target
After=systemd-user-sessions.service

[Service]
Type=simple
ExecStart=/opt/bin/anydesk --service
PIDFile=/run/anydesk.pid
KillMode=mixed
TimeoutStopSec=30
User=root
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
