# sys-fs/zfs/files/zed
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need zfs
}

start() {
	ebegin "Starting ZFS Event daemon"

	checkpath -q -d -m 0755 /var/run/zed

	start-stop-daemon --start -q \
		--exec /sbin/zed -- -M \
			-p /var/run/zed/pid
	eend $?

}

stop() {
ebegin "Stopping ZFS Event daemon"
	start-stop-daemon --stop -q --pidfile /var/run/zed/pid
eend $?
}
