# sys-fs/zfs/files/zed
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-fs/zfs/files/zed,v 1.1 2015/04/27 13:58:25 ryao Exp $

depend() {
	need zfs
}

start() {
	ebegin "Starting ZFS Event daemon"

	checkpath -q -d -m 0755 /var/run/zed

	start-stop-daemon --start -q \
		--exec /sbin/zed -- -M \
			-p /var/run/zed/pid
	eend $?

}

stop() {
ebegin "Stopping ZFS Event daemon"
	start-stop-daemon --stop -q --pidfile /var/run/zed/pid
eend $?
}
