# net-wireless/bluez/files/bluetooth-init.d-r2
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	after coldplug
	need dbus localmount hostname
}

start() {
   	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --exec /usr/sbin/bluetoothd
	eend $?
}

stop() {
	ebegin "Shutting down ${SVCNAME}"
	start-stop-daemon --stop --quiet --exec /usr/sbin/bluetoothd
	eend $?
}
