# net-proxy/squid/files/squid.cron
#!/bin/sh
if test -e /var/run/squid.pid ; then
	test -n "$(cat /var/run/squid.pid|xargs ps -p|grep squid)" && \
		/usr/sbin/squid -k rotate
fi
