# app-misc/vcontrold/files/vcontrold
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/vcontrold"
