# dev-libs/libcgroup/files/cgred.initd-r1
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="/usr/sbin/cgrulesengd"
command_args="${CGRED_OPTS}"
command_background="1"
pidfile="/var/run/cgred.pid"
description="CGroups Rules Engine Daemon"
extra_started_commands="reload"

depend()
{
	need cgconfig
	use logger
}

reload()
{
	ebegin "Reloading CGroup Rules Engine Daemon"
	kill -USR2 $(cat "${pidfile}")
	eend $?
}
