# net-analyzer/zabbix/files/1.6.6/init.d/zabbix-proxy
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
	#use mysql postgresql
}

start() {
	ebegin "Starting Zabbix proxy"
	start-stop-daemon --start --user zabbix --group zabbix --exec /usr/sbin/zabbix_proxy
	eend $?
}

stop() {
	ebegin "Stopping Zabbix proxy"
	start-stop-daemon --stop --pidfile /var/run/zabbix/zabbix_proxy.pid
	eend $?
}
