# sys-process/incron/files/incrond.service
[Unit]
Description=Inotify System Scheduler

[Service]
ExecStart=/usr/sbin/incrond --foreground -f /etc/incron.conf

[Install]
WantedBy=multi-user.target

