# net-misc/igmpproxy/files/igmpproxy.service
[Unit]
Description=IGMP proxy
After=network.target

[Service]
ExecStart=/usr/sbin/igmpproxy /etc/igmpproxy.conf

[Install]
WantedBy=multi-user.target
