# sys-apps/nca/files/ncad.initd
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --quiet --exec /sbin/ncad
	eend ${?}
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --retry 5 --quiet --exec /sbin/ncad
	eend ${?}
}
