# net-misc/dibbler/files/dibbler-client.service
[Unit]
Description=Dibbler IPv6 DHCPv6 Client

[Service]
ExecStart=/usr/sbin/dibbler-client run

[Install]
WantedBy=multi-user.target
