# net-misc/dibbler/files/dibbler-relay
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting DHCPv6 relay: dibbler-relay"
	/usr/sbin/dibbler-relay start
	eend $?
}

stop() {
	ebegin "Stopping DHCPv6 relay: dibbler-relay"
	/usr/sbin/dibbler-relay stop
	eend $?
}
