# app-misc/g15daemon/files/99-g15daemon.rules
SUBSYSTEM=="usb", ATTR{product}=="G15 Gaming Keyboard", RUN+="/usr/bin/g15daemon-hotplug"
