# sys-process/anacron/files/anacron.rc6
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	use logger
	need clock hostname 
}

start() {
	ebegin "Running anacron"
	/usr/sbin/anacron -s >>/var/log/cron.log 2>&1
	eend $?
}
