# sys-cluster/glusterfs/files/glusterd-r1.initd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="Gluster elastic volume management daemon"
command="/usr/sbin/glusterd"
pidfile="/var/run/${SVCNAME}.pid"
command_args="-N"

command_background="yes"

depend() {
	need net
	before netmount
}
