# app-admin/gkrellm/files/gkrellmd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
	after lm_sensors
	after hddtemp
}

start() {
	ebegin "Starting GNU Krell Monitor daemon"
	start-stop-daemon --start --quiet --background --pidfile /var/run/gkrellmd.pid --make-pidfile --exec /usr/bin/gkrellmd -- ${GKRELLMD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping GNU Krell Monitor daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/gkrellmd.pid --name gkrellmd
	eend $?
}
