# app-admin/gkrellm/files/gkrellmd.service
[Unit]
Description=GNU Krell Monitors server
Documentation=man:gkrellmd(1)
After=network.target

[Service]
ExecStart=/usr/bin/gkrellmd

[Install]
WantedBy=multi-user.target
