# sys-libs/libtrash/files/cleanTrash.cron
#!/bin/sh
#
# Uncomment the following if you want
# your Trash cleaned out regularly.
#
# This is part of the libtrash package.
#

#if [ -x /usr/sbin/ct2.pl ]
#then
#	/usr/sbin/ct2.pl > /dev/null
#fi
