# sys-apps/rsbac-admin/files/rklogd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/rsbac-admin/files/rklogd,v 1.3 2012/05/17 23:17:00 blueness Exp $

daemon="rklogd"
rsbac_info="/proc/rsbac-info"

start() {
	ebegin "Starting rklogd"

	local ret
	if [ -d "${rsbac_info}" ]; then
		if type run-jail >/dev/null 2>&1; then
			$(which run-jail) $daemon start-stop-daemon --start --exec $(which ${daemon})
		else
			start-stop-daemon --start --exec $(which ${daemon})
		fi
		ret=$?
	else
		ewarn "No rsbac-info directory avaible."
		ret=1 
	fi

	eend $ret
}

stop() {
	ebegin "Stopping rklogd"
	start-stop-daemon --stop --exec $(which ${daemon})
	eend $?
}
