# net-misc/pingu/files/pingu.service
[Unit]
Description="Policy routing daemon"
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/pingu -d
Restart=on-abort

[Install]
WantedBy=multi-user.target
