# media-radio/direwolf/files/direwolf.service
[Unit]
Description=Dire Wolf

[Service]
Type=simple
ExecStart=/usr/bin/direwolf -l /var/log/direwolf
Restart=always

[Install]
WantedBy=multi-user.target
