# sys-auth/consolekit/files/consolekit-1.0.0.initd
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

description="consolekit - tracks login sessions and seats"
pidfile="/run/ConsoleKit/pid"
command="/usr/sbin/console-kit-daemon"

depend() {
	need dbus
	use logger
}

start_pre() {
	checkpath -q -d -m 0755 /run/ConsoleKit
}
