# media-plugins/vdr-xineliboutput/files/vdr-frontend
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="Start xineliboutput sxfe/sxfb remote frontend"

source /etc/conf.d/vdr.xineliboutput

command="/usr/bin/vdr-sxfe"
command_args="${REMOTE_FRONTEND}"
pidfile="/var/run/vdr-frontend.pid"
start_stop_daemon_args="--user vdr --group vdr --make-pidfile --background"
