# dev-util/buildbot-worker/files/buildbot_worker.target
[Unit]
Description=Buildbot-Worker target that allows to start or stop all buildbot_worker@*.service instances

[Install]
WantedBy=multi-user.target
