# net-dns/dnsmasq/files/dnsmasq-init-r5
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

extra_started_commands="reload rotate"

# Convert to uppercase (e.g. DNSMASQ.ETH1)
PREFIX="${RC_SVCNAME^^}"
# Replace dots with underscores {e.g. DNSMASQ_ETH1}
PREFIX="${RC_SVCNAME//./_}"

# Construct variable names
INSTANCE_USER="${PREFIX}_USER"
INSTANCE_GROUP="${PREFIX}_GROUP"
INSTANCE_OPTS="${PREFIX}_OPTS"

# Dereference needed variables (default user:group to dnsmasq:dnsmasq)
DNSMASQ_USER="${!INSTANCE_USER-dnsmasq}"
DNSMASQ_GROUP="${!INSTANCE_GROUP-dnsmasq}"
DNSMASQ_OPTS="${!INSTANCE_OPTS}"

# Prepend config file for non-default instance
# This way it can be overriden in DNSMASQ_OPTS
INSTANCE="${RC_SVCNAME#*.}"
if [[ -n "${INSTANCE}" ]]
then
	DNSMASQ_OPTS="--conf-file=/etc/${RC_SVCNAME}.conf ${DNSMASQ_OPTS}"
fi

pidfile="/var/run/${RC_SVCNAME}.pid"
command="/usr/sbin/dnsmasq"
command_args="-x ${pidfile} ${DNSMASQ_OPTS}"
retry="TERM/3/TERM/5"

depend() {
	provide dns
	need localmount net
	after bootmisc
	use logger
}

reload() {
	ebegin "Reloading ${RC_SVCNAME}"
	start-stop-daemon --signal HUP --pidfile "${pidfile}"
	eend $?
}

rotate() {
	ebegin "Reopening ${RC_SVCNAME} log file"
	start-stop-daemon --signal USR2 --pidfile "${pidfile}"
	eend $?
}
