# net-voip/gnugk/files/gnugk.rc6
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting H.323 Gatekeeper"
	start-stop-daemon --start --quiet --background \
		--exec /usr/sbin/gnugk -- ${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping H.323 Gatekeeper"
	start-stop-daemon --stop --quiet --pidfile /var/run/gnugk.pid
	eend $?
}
