# app-backup/bacula/files/newscripts/bacula-fd.initd
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command=/usr/sbin/bacula-fd
command_args="${FD_OPTIONS}"
# make sure the following line matches the portnumber from bacula-fd.conf
pidfile=/var/run/bacula-fd.9102.pid

depend() {
	need net
	use dns
}
