# www-apps/ttyd/files/ttyd.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

name="ttyd daemon"
pidfile="/var/run/ttyd.pid"
command="/usr/bin/ttyd"
command_args="${TTYD_OPTS}"
command_background="true"
start_stop_daemon_args="-b"

depend() {
   use logger
}
