# sys-power/athcool/files/athcool
#!/sbin/runscript
# $Id$

start() {
	ebegin "Starting ${SVCNAME}"
	athcool on >/dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	athcool off >/dev/null 2>&1
	eend $?
}

