# sys-power/athcool/files/athcool
#!/sbin/runscript
# $Header: /var/cvsroot/gentoo-x86/sys-power/athcool/files/athcool,v 1.2 2007/10/01 10:57:44 uberlord Exp $

start() {
	ebegin "Starting ${SVCNAME}"
	athcool on >/dev/null 2>&1
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	athcool off >/dev/null 2>&1
	eend $?
}

