# net-misc/stuntman/files/initd.v1
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting STUNTMAN"
	start-stop-daemon --start --quiet --exec /usr/sbin/stunserver \
		--make-pidfile --pidfile=/run/stunserver.pid --background
	eend $?
}

stop() {
	ebegin "Stopping STUNTMAN"
	start-stop-daemon --stop --quiet --pidfile=/run/stunserver.pid
	eend $?
}
