# x11-apps/radeon-profile-daemon/files/radeon-profile-daemon.initd-r2
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Daemon for radeon-profile GUI"

command="/usr/sbin/radeon-profile-daemon"
command_background="true"
pidfile="/run/${SVCNAME}/${SVCNAME}.pid"

radeon_socket="/run/${SVCNAME}/${SVCNAME}-server"

start_pre() {
	checkpath -q -d -o root:video -m 0750 /run/${SVCNAME}
}
