# sys-devel/prelink/files/prelink.timer
[Unit]
Description=Daily system-wide prelinking

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=timers.target
