# sys-devel/prelink/files/prelink.service
[Unit]
Description=Perform system-wide prelinking

[Service]
Type=oneshot
ExecStart=/etc/cron.daily/prelink
