# mail-mta/exim/files/exim_at.service
[Unit]
Description=Exim Mail Daemon per-connection server

[Service]
ExecStart=-/usr/bin/exim -bs
StandardInput=socket
StandardError=syslog
