# mail-mta/exim/files/exim-submission.socket
[Unit]
Description=Exim Mail Transfer Agent (message submission)
Conflicts=exim.service

[Socket]
ListenStream=587
Accept=yes

[Install]
WantedBy=sockets.target
