# net-misc/ktls-utils/files/tlshd.initd
#!/sbin/openrc-run
# Copyright 2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="tlshd daemon"
description="Handshake service for kernel TLS consumers"
command="/usr/sbin/tlshd"
pidfile="/run/tlshd.pid"

command_background="yes"

depend() {
    need net
}
