# sci-electronics/labone/files/labone-data-server
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile="/run/labone-data-server.pid"
command="/opt/bin/ziDataServer"
command_args=""

depend() {
    need network
}
