# net-misc/lldpd/files/lldpd-initd-5
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

name=lldpd
pidfile=/run/lldpd.pid
command=/usr/sbin/lldpd
retry="TERM/10/KILL/5"
command_args="${LLDPD_OPTS}"
start_stop_daemon_args="--wait 300"

depend() {
	use net
}

start_pre() {
	checkpath --directory /run/lldpd -o lldpd -m 0700
}

# vim:ft=gentoo-init-d:noet:ts=4:sts=4:sw=4:
