# net-irc/shadowircd/files/shadowircd.confd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/shadowircd/files/shadowircd.confd,v 1.2 2012/03/03 07:10:09 jdhore Exp $

# Specify any extra options to pass to shadowircd
SHADOWIRCD_OPTS=""
