# www-apps/sonarr/files/sonarr.init-r1
#!/sbin/openrc-run

pidfile="/run/sonarr.pid"
command="/opt/sonarr/Sonarr"
command_user="sonarr:sonarr"
command_background=true

