# mail-mta/opensmtpd/files/smtpd.initd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="smtpd"
pidfile="/run/smtpd.pid"
description="SMTP daemon from OpenBSD"

depend() {
	need net
}
