# mail-mta/opensmtpd/files/smtpd.socket
[Unit]
Description=OpenSMTPD Socket

[Socket]
ListenStream=/var/run/smtpd.sock

[Install]
WantedBy=sockets.target
