# net-nds/gssproxy/files/gssproxy
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="gssproxy"
command_args="-D"
pidfile="/var/run/gssproxy.pid"

depend() {
	use logger
}

start_pre() {
	checkpath -d -m 0755 /var/lib/gssproxy
}
