# net-irc/charybdis/files/charybdis.confd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/charybdis/files/charybdis.confd,v 1.1 2012/03/03 06:14:56 binki Exp $

# Specify any extra options to pass to charybdis
CHARYBDIS_OPTS=""
