# www-apps/bugzilla/files/bugzilla.cron.daily
#!/bin/sh
cd /var/www/bugzilla/
./collectstats.pl
./whineatnews.pl
echo "*** Bugzilla nighlty flush has done! ***"
