# games-simulation/openttd/files/openttd.initd-r2
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

description="OpenTTD dedicated game server"
pidfile="/run/${SVCNAME}.pid"

command="/usr/bin/openttd"
command_args="-D ${openttd_opts}"
command_background="true"
command_group="openttd"
command_user="openttd"

depend() {
	need net
}
