# net-p2p/rtorrent/files/rtorrent-r1.init
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2

description="rTorrent BitTorrent client"

command="/usr/bin/rtorrent"
command_args="-o system.daemon.set=true"
command_background=true
command_user="${USER}"
pidfile="/run/rtorrent.pid"

depend()
{
   use net ypbind nis
   after slapd mysqld postgresql
}
