# net-p2p/rtorrent/files/rtorrentd_at.service
[Unit]
Description=rTorrent
Requires=network.target local-fs.target

[Service]
Type=forking
RemainAfterExit=yes
KillMode=none
User=%I
ExecStart=/usr/bin/screen -d -m -S rtorrentd /usr/bin/rtorrent
ExecStop=/usr/bin/screen -S rtorrentd -X quit
WorkingDirectory=/home/%I/

[Install]
WantedBy=multi-user.target
