# app-backup/bareos/files/bareos-fd.initd
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	use dns
}

command="/usr/sbin/bareos-fd"
command_args="${FD_OPTIONS}"
pidfile="/run/bareos/bareos-fd.9102.pid"
