# media-video/coherence/files/coherence-init
#!/sbin/runscript

depend() {
  need net
  after bootmisc
}

start() {
  ebegin "Starting Coherence"
  start-stop-daemon --start --exec /usr/bin/coherence \
    --pidfile /var/run/coherence.pid \
    --background \
    -- -c /etc/coherence.conf
  eend $?
}

stop() {
  ebegin "Stopping Coherence"
  start-stop-daemon --stop --exec /usr/bin/coherence \
    --pidfile /var/run/coherence.pid
  eend $?
}

