# net-misc/etherdfs/files/etherdfs.initd
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="etherdfs daemon"
command="/usr/bin/ethersrv-linux"
command_args="${INTERFACE} ${STORAGEDIR}"

depend() {
	need net
}
