# net-analyzer/sinfo/files/sinfod.initd
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="/usr/sbin/${SVCNAME}"
command_args="${SINFOD_OPTS}"

depend() {
	need net
}
