# net-nds/ypserv/files/rpc.ypxfrd-2.23
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net portmap ypserv
}

start() {
	ebegin "Starting rpc.ypxfrd"
	YOPTS=""
	if [ -n "${YPXFRD_OPTS}" ]; then YOPTS="-- ${YPXFRD_OPTS}"; fi
 	start-stop-daemon --start --quiet --pidfile /var/run/ypxfrd.pid --exec /usr/sbin/rpc.ypxfrd ${YOPTS}
	eend $?
}

stop() {
	ebegin "Stopping rpc.ypxfrd"
	start-stop-daemon --stop --quiet --pidfile /var/run/ypxfrd.pid --exec /usr/sbin/rpc.ypxfrd
	eend $?
}
