# net-nds/ypserv/files/ypserv
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net portmap
	use domainname
}

start() {
	ebegin "Starting NIS Server"
	YOPTS=""
	if [ -n "${YPSERV_OPTS}" ]; then YOPTS="-- ${YPSERV_OPTS}"; fi
	start-stop-daemon --start --quiet --exec /usr/sbin/ypserv ${YOPTS}
	eend $?
}

stop() {
	ebegin "Stopping NIS Server"
	start-stop-daemon --stop --quiet --exec /usr/sbin/ypserv
	eend $?
}
