# dev-python/flower/files/flower.service
[Unit]
Description=Real-time monitor and web admin for Celery distributed task queue
Requires=network-online.target
After=network-online.target

[Service]
User=flower
PrivateDevices=Yes
Environment=_FLOWER_OPTS="--conf=/etc/flower/config.py"
ExecStart=/usr/bin/flower $_FLOWER_OPTS

[Install]
WantedBy=multi-user.target
