# app-containers/incus/files/incus-0.1.socket
[Unit]
Description=incus - unix socket
Documentation=man:incus(1)

[Socket]
ListenStream=/var/lib/incus/unix.socket
SocketGroup=incus
SocketMode=0660
Service=incus.service

[Install]
WantedBy=sockets.target
