# net-ftp/tftp-hpa/files/in.tftpd.rc6
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-ftp/tftp-hpa/files/in.tftpd.rc6,v 1.2 2005/07/30 06:29:14 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting tftpd"
	/usr/sbin/in.tftpd -l ${INTFTPD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping tftpd"
	start-stop-daemon --stop --exec /usr/sbin/in.tftpd
	eend $?
}
