# net-ftp/tftp-hpa/files/tftp.service
[Unit]
Description=Tftp Server

[Service]
ExecStart=/usr/sbin/in.tftpd -R 4096:32767 -s /tftproot
StandardInput=socket
