# dev-db/hsqldb/files/hsqldb
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
    use net
}

start() {
	ebegin "Starting HSQL Database"
	hsqldb_enable=yes /var/lib/hsqldb/bin/hsqldb start
	eend $?
}

stop() {
	ebegin "Stopping HSQL Database"
	hsqldb_enable=yes /var/lib/hsqldb/bin/hsqldb stop
	eend $?
}
