# net-wireless/pi-bluetooth/files/init.d_rpi-hciuart-1
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Raspberry Pi hciuart activation service"
command="/usr/bin/btuart"

depend() {
	need localmount
	before bluetooth
}

stop() {
        # Nothing to stop
        return 0;
}
