# net-dns/maradns/files/maradns.service.in
[Unit]
Description=maradns DNS server
After=network.target

[Service]
ExecStart=@bindir@/maradns -cf @sysconfdir@/mararc

[Install]
WantedBy=multi-user.target
