# net-dns/maradns/files/zoneserver.rc6
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting zoneserver"
	start-stop-daemon --start --quiet --name zoneserver \
		--exec /usr/bin/duende /usr/sbin/zoneserver
	eend $?
}

stop() {
	ebegin "Stopping zoneserver"
	start-stop-daemon --stop --quiet --name zoneserver
	eend $?
}
