# media-sound/rescrobbled/files/rescrobbled.initd
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="rescrobbled daemon"
description="MPRIS music scrobbler daemon"
command=/usr/bin/rescrobbled
procname="rescrobbled"
