# dev-db/keydb/files/keydb.initd-6
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

: ${KEYDB_DIR:=/var/lib/keydb}
: ${KEYDB_CONF:=/etc/keydb/keydb.conf}
: ${KEYDB_OPTS:="${KEYDB_CONF}"}
: ${KEYDB_USER:=keydb}
: ${KEYDB_GROUP:=keydb}
: ${KEYDB_TIMEOUT:=30}

# https://bugs.gentoo.org/631002#c10
# Force '--daemonize no' to override the config file
command="/usr/sbin/keydb-server"
command_args="${KEYDB_OPTS} --daemonize no"
command_background="true"
command_user="${KEYDB_USER}:${KEYDB_GROUP}"
pidfile="/run/${RC_SVCNAME}.pid"
retry="${KEYDB_TIMEOUT}"
start_stop_daemon_args="--chdir \"${KEYDB_DIR}\""

depend() {
	use localmount logger
	after keepalived
}
