# media-video/flumotion/files/flumotion-init-0.2.0
#!/sbin/runscript

depend() {
    need net
}

start() {
    ebegin "Starting Flumotion"
    HOME=/usr/share/flumotion start-stop-daemon -o --start --chuid flumotion:flumotion --exec /usr/sbin/flumotion start 
    eend $?
}

stop() {
    ebegin "Stopping Flumotion"
    /usr/sbin/flumotion stop
    eend $?
}
