# app-admin/verynice/files/verynice.service
[Unit]
Description=A tool for dynamically adjusting the nice-level of processes

[Service]
ExecStart=/usr/sbin/verynice

[Install]
WantedBy=multi-user.target
