# app-admin/verynice/files/verynice
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/verynice/files/verynice,v 1.4 2008/02/01 12:31:24 flameeyes Exp $

depend() {
	after checkroot root
}

start() {
	ebegin "Starting verynice"
	/usr/sbin/verynice -d /var/run/verynice.pid
	eend $?
}

stop() {
	ebegin "Stopping verynice"
	start-stop-daemon --stop --quiet --pidfile=/var/run/verynice.pid
	eend $?
}
