# www-apps/sonarr-bin/files/sonarr.init-r2
#!/sbin/openrc-run

pidfile="/run/sonarr.pid"
command="/opt/sonarr/Sonarr"
command_user="${GROUP:-sonarr}:${USER:-sonarr}"
command_background=true

