# sys-fs/iprutils/files/iprdump
#!/sbin/openrc-run

depend() {
	need logger iprinit
}

start() {
	ebegin "Starting the ipr dump daemon"
	start-stop-daemon --start  --quiet --exec /sbin/iprdump

}

stop() {
	ebegin "Stopping the ipr dump daemon"
	start-stop-daemon --stop --name iprdump
}
