# sys-fs/iprutils/files/iprupdate
#!/sbin/openrc-run

depend() {
	need logger iprinit
}

start() {
	ebegin "Checking ipr microcode levels"
	/sbin/iprupdate 2>&1 > /dev/null
	eend $? "Failed to complete ipr microcode updates"
}
