# sys-apps/earlyoom/files/earlyoom-r1
#!/sbin/openrc-run
name=earlyoom
command=/usr/bin/earlyoom
command_args="${EARLYOOM_ARGS}"
command_background=1
pidfile=/run/earlyoom.pid
