# net-irc/ultimate/files/ultimateircd.rc
#!/sbin/openrc-run
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	use dns net
}

start() {
	ebegin "Starting ultimate-ircd"
	start-stop-daemon --start --quiet --exec /usr/bin/ultimate-ircd \
		--chuid ${ULTIMATEIRCD_USER} &>/dev/null
	eend $?
}

stop() {
	ebegin "Shutting down ultimate-ircd"
	killall ultimate-ircd
	eend $?
}
