# sys-fs/multipath-tools/files/rc-multipathd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need localmount
	after modules
}

start() {
	ebegin "Starting multipathd"
	start-stop-daemon --start --quiet --exec /sbin/multipathd
	eend $?
}

stop() {
	ebegin "Stopping multipathd"
	start-stop-daemon --stop --quiet --pidfile /var/run/multipathd.pid
	eend $?
}
