# net-analyzer/munin/files/munin-graph.service
[Unit]
Description=Munin FastCGI Graph
Requires=munin-graph.socket

[Service]
User=munin
Group=munin
ExecStart=/usr/libexec/munin/munin-graph
StandardInput=socket
StandardOutput=inherit
StandardError=journal
