# net-analyzer/munin/files/munin-html.socket
[Unit]
Description=Munin FastCGI HTML Socket

[Socket]
ListenStream=/run/munin/fcgi-html.sock

[Install]
WantedBy=sockets.target
