# net-misc/openvswitch/files/ovs-controller-r1
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="Open vSwitch OpenFlow controller"

command="/usr/bin/ovs-controller"
command_args="
    --pidfile
    --detach
    --monitor
    ${OPTIONS} ${METHODS}"
pidfile="/var/run/openvswitch/ovs-controller.pid"

depend() {
    need net
    use logger
}

start_pre() {
	checkpath -d "/var/run/openvswitch" -m 0750
}
