# net-misc/openvswitch/files/ovs-controller.service
[Unit]
Description=Open vSwitch OpenFlow controller
Requires=network.target

[Service]
ExecStart=/usr/bin/ovs-controller --monitor ptcp:

[Install]
WantedBy=multi-user.target
