# net-misc/openvswitch/files/ovs-vswitchd-r1
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="Open vSwitch virtual switch"

command="/usr/sbin/ovs-vswitchd"
command_args="
    --pidfile
    --detach
    --monitor
    ${OPTIONS} ${DATABASE}"
pidfile="/var/run/openvswitch/ovs-vswitchd.pid"

depend() {
	use ovsdb-server logger
}

start_pre() {
	checkpath -d "/var/run/openvswitch" -m 0750
}
