# sys-apps/timer_entropyd/files/timer_entropyd.initd.1
#!/sbin/openrc-run
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command=/usr/libexec/timer_entropyd

depend() {
	need localmount
	before net
	use logger

	provide entropy
}
