# app-laptop/i8kutils/files/i8kmon.service
[Unit]
Description=Dell laptop thermal monitoring
ConditionPathExists=/proc/i8k

[Service]
ExecStart=/usr/bin/i8kmon -d

[Install]
WantedBy=multi-user.target
