# net-analyzer/net-snmp/files/snmpd.service
[Unit]
Description=Simple Network Management Protocol (SNMP) Daemon
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/snmpd -f

[Install]
WantedBy=multi-user.target
