# sys-devel/distcc/files/2.18/init
#!/sbin/runscript
# $Header: /var/cvsroot/gentoo-x86/sys-devel/distcc/files/2.18/init,v 1.1 2004/10/12 17:21:43 lisa Exp $

depend() {
	need net
	use ypbind
}

start() {
	[ -e "${DISTCCD_PIDFILE}" ] && rm -f ${DISTCCD_PIDFILE} &>/dev/null
	
	ebegin "Starting distccd"
	chown distcc `dirname ${DISTCCD_PIDFILE}` &>/dev/null
	TMPDIR="${TMPDIR}" \
	PATH="$(gcc-config --get-bin-path):${PATH}" \
	/sbin/start-stop-daemon --start --quiet --startas ${DISTCCD_EXEC} \
	--pidfile ${DISTCCD_PIDFILE} -- \
	--pid-file ${DISTCCD_PIDFILE} -N ${DISTCCD_NICE} --user distcc \
	${DISTCCD_OPTS}

	eend $?
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile "${DISTCCD_PIDFILE}"
	rm -f "${DISTCCD_PIDFILE}"
	eend $?
}

