# app-misc/lcdproc/files/0.5.2-r2-lcdproc.initd
#!/sbin/openrc-run
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend(){
	use LCDd
	after ntp-client
}

start() {
	ebegin "Starting lcdproc"
	start-stop-daemon --start --background \
		--exec /usr/bin/lcdproc
	eend $?
}

stop() {
	ebegin "Stopping lcdproc"
	start-stop-daemon --stop --exec /usr/bin/lcdproc
	eend $?
}
