# app-misc/lcdproc/files/lcdproc.service
[Unit]
Description=LCDProc (LCDd)
After=network-online.target

[Service]
ExecStart=/usr/bin/lcdproc -c /etc/lcdproc.conf
Type=simple

[Install]
WantedBy=multi-user.target
