# sys-apps/irqbalance/files/irqbalance.init.2
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Id$

depend() {
	need localmount
	after bootmisc
}

command="/usr/sbin/irqbalance"
command_args="${IRQBALANCE_OPTS}"
