# sys-cluster/ganglia/files/gmetad.rc
#!/sbin/openrc-run

depend() {
	need net
	after ntp-client ntpd
}

start() {
	ebegin "Starting GANGLIA gmetad: "
	start-stop-daemon --start --quiet --exec /usr/sbin/gmetad
	eend $? "Failed to start gmetad"
}

stop() {
	ebegin "Shutting down GANGLIA gmetad: "
	start-stop-daemon --stop --quiet --exec /usr/sbin/gmetad
	eend $? "Failed to stop gmetad"
}
