# app-laptop/hdapsd/files/hdapsd.init.3
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2

depend() {
    need localmount
}

checkconfig() {
    if [ -n "${DISK}" ]; then
	ewarn "Please migrate DISK to DISKLIST in /etc/conf.d/hdapsd."
	DISKLIST=${DISK}
    fi

    if [ -z "${THRESHOLD}" ] ; then
	eerror "You should setup THRESHOLD in /etc/conf.d/hdapsd."
	return 1
    fi

    #empty DISKLIST is allowed as well
    for DISK in ${DISKLIST}; do
        if [ ! -b /dev/${DISK} ]; then
            eerror "Could not find disk /dev/${DISK}!"
            eerror "Adjust the DISK setting in /etc/conf.d/hdapsd"
            return 1
        fi

        if [ ! -e /sys/block/${DISK}/queue/protect ] && [ ! -e /sys/block/${DISK}/device/unload_heads ] ; then
            eerror "No protect entry for ${DISK}!"
            eerror "Make sure your kernel is patched with the blk_freeze patch"
            return 1
        fi
    done

    # Load the tp_smapi module first
    # This is not a requirement, but it helps hdapsd adaptive mode
    if [ ! -e /sys/devices/platform/smapi ] ; then
	modprobe tp_smapi 2>/dev/null
    fi
	
    if [ ! -d /sys/devices/platform/hdaps ] ; then
	ebegin "Loading hdaps module"
        modprobe hdaps
	eend $? || return 1
    fi
}

start() {
    checkconfig || return 1

    local DISKOPTS DISK
    for DISK in ${DISKLIST}; do
	DISKOPTS="${DISKOPTS} -d ${DISK}"
    done

    ebegin "Starting Hard Drive Active Protection System daemon"
    start-stop-daemon --start --exec /usr/sbin/hdapsd \
	--pidfile /var/run/hdapsd.pid \
	-- --syslog --background --pidfile \
	${DISKOPTS} --sensitivity "${THRESHOLD}" ${OPTIONS}
    eend $?
}

stop() {
    ebegin "Stopping Hard Drive Active Protection System daemon"
    start-stop-daemon --stop --exec /usr/sbin/hdapsd \
	--pidfile /var/run/hdapsd.pid
    eend $?
}
