# x11-misc/dunst/files/dunst.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need dbus
}

export DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"
supervisor=supervise-daemon
command="/usr/bin/dunst"
