# mail-mta/nullmailer/files/nullmailer.service
[Unit]
Description=Nullmailer relay-only MTA
Requires=network.target
After=local-fs.target syslog.target
ConditionPathExists=/var/nullmailer/queue
ConditionPathExists=/etc/nullmailer/defaultdomain
ConditionPathExists=/etc/nullmailer/me

[Service]
WorkingDirectory=/var/nullmailer
ExecStart=/usr/sbin/nullmailer-send -s
User=nullmail
Group=nullmail
Restart=always

[Install]
WantedBy=multi-user.target
