# sys-apps/pcsc-tools/files/smartcard.cron
#!/bin/sh

set -e

SRC=http://ludovic.rousseau.free.fr/softwares/pcsc-tools/smartcard_list.txt
DEST=/usr/share/pcsc/smartcard_list.txt

touch ${DEST} >/dev/null 2>&1 || exit 1

if ! wget --quiet --no-timestamping -O "${DEST}".new "$SRC"; then
	rm -f "${DEST}".new
fi

chmod a+r "${DEST}".new
mv "${DEST}".new "${DEST}"
