# net-misc/netkit-rwho/files/netkit-rwho-0.17-rc
#!/sbin/openrc-run

depend() {
	need net
}

start() {
	ebegin "Starting rwhod"
	/sbin/start-stop-daemon --start --quiet --exec /usr/sbin/rwhod -- $RWHOD_OPTIONS
	eend $?
}

stop() {
	ebegin "Stopping rwhod"
	/sbin/start-stop-daemon --stop --quiet --name rwhod
	eend $?
}
