# app-i18n/skkserv/files/skkserv.initd
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

start() {
	ebegin "Starting skkserv"
	start-stop-daemon --start --quiet --exec /usr/sbin/skkserv
	eend $?
}

stop() {
	ebegin "Stopping skkserv"
	start-stop-daemon --stop --quiet --exec /usr/sbin/skkserv
	eend $?
}
