# net-misc/portfwd/files/portfwd.service
[Unit]
Description=Port Forwarding Daemon

[Service]
ExecStart=/usr/sbin/portfwd -g -f -t

[Install]
WantedBy=multi-user.target
