# dev-util/jenkins-bin/files/jenkins-bin-r2.service
[Unit]
Description=Jenkins Daemon

[Service]
EnvironmentFile=/etc/conf.d/jenkins
User=jenkins
Type=simple
ExecStart=/bin/sh -c 'eval ${JENKINS_COMMAND_LINE}'
SyslogIdentifier=jenkins
#MemoryLimit=768M
OOMScoreAdjust=-1000
LimitCPU=infinity
LimitFSIZE=infinity
LimitDATA=infinity
LimitCORE=0
LimitAS=infinity
LimitLOCKS=infinity

[Install]
WantedBy=multi-user.target
