# net-vpn/fp-multiuser/files/fp-multiuser.service
[Unit]
Description=A frp server plugin to support multiple users for frp
After=network.target

[Service]
User=fp-multiuser
Group=fp-multiuser
Environment=FP_MULTIUSER_OPTS="-f /etc/fp-multiuser/tokens -l 127.0.0.1:7200"
Restart=on-failure
RestartSec=5s
ExecStart=/usr/bin/fp-multiuser $FP_MULTIUSER_OPTS

[Install]
WantedBy=multi-user.target
