# net-misc/openssh/files/ssh-agent.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

supervisor=supervise-daemon
command="/usr/bin/ssh-agent"
command_args="-D -a ${XDG_RUNTIME_DIR}/ssh-agent.sock"
