# net-wireless/eiwd/files/iwd.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/libexec/iwd"
pidfile="/run/iwd.pid"
command_background="yes"

depend() {
	provide net
	after bootmisc modules
	before dns dhcpcd
	keyword -shutdown
}
