# sys-apps/busybox/files/ntpd.initd
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="/bin/busybox ntpd"
command_args="${NTPD_OPTS}"
pidfile="/var/run/ntpd.pid"

depend() {
	use net dns logger
	after ntp-client
}
