# sys-apps/busybox/files/syslogd.initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/busybox/files/syslogd.initd,v 1.2 2012/12/21 03:15:28 blueness Exp $

command="/bin/busybox syslogd"
command_args="${SYSLOGD_OPTS}"
pidfile="/var/run/syslogd.pid"

depend() {
	need clock hostname
	provide logger
}
