# sys-auth/seatd/files/seatd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon
command="seatd"
command_args="-g video"
