# net-firewall/firewalld/files/firewalld.init
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="FirewallD"
command=/usr/sbin/firewalld
pidfile=/var/run/firewalld.pid

depend() {
	need dbus
	provide iptables ip6tables ebtables
}
