# app-editors/emacs/files/emacs.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

supervisor=supervise-daemon
description="Emacs text editor"

command="@EMACS_BIN@"
command_args_foreground="--fg-daemon"
