# sys-power/power-profiles-daemon/files/power-profiles-daemon.initd
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="power-profiles-daemon"
description="Makes power profiles handling available over D-Bus"

pidfile="/run/power-profiles-daemon.pid"
command="/usr/libexec/power-profiles-daemon"
command_background=true

start_pre() {
	checkpath -d /var/lib/power-profiles-daemon
}
