# app-containers/lxd/files/lxd-containers-4.0.0.service
[Unit]
Description=LXD - container startup/shutdown
Documentation=man:lxd(1)
After=lxd.socket lxd.service
Requires=lxd.socket

[Service]
Type=oneshot
ExecStart=/usr/sbin/lxd activateifneeded
ExecStop=/usr/sbin/lxd shutdown
TimeoutStartSec=600s
TimeoutStopSec=600s
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
