# media-sound/sndio/files/sndiod
#!/sbin/openrc-run
command="sndiod"
command_user="sndiod"
