# app-backup/amanda/files/amanda.socket
[Unit]
Description=Amanda Socket
[Socket]
ListenStream=10080
Accept=true
[Install]
WantedBy=sockets.target
