# games-fps/ut2004/files/ut2004-ded.initd
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile="/var/run/${RC_SVCNAME}.pid"
command_background="true"
command_user="${ut2004_ded_user}:${ut2004_ded_group}"
command="ut2004-ded"
command_args="${ut2004_ded_opts}"

depend() {
	need net
}
