# www-apps/prowlarr/files/prowlarr.init
#!/sbin/openrc-run

pidfile="/run/prowlarr.pid"
command="/opt/prowlarr/Prowlarr"
command_user="prowlarr:prowlarr"
command_background=true

