# net-dns/dnrd/files/dnrd
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Id$

PIDFILE="/var/run/dnrd.pid"

depend() {
	provide dns
	need net
}

start() {
	ebegin "Starting dnrd"
	/usr/sbin/dnrd $DNRD_OPTS &> /dev/null &
	eend $?
}

stop() {
	ebegin "Stopping dnrd"
	/usr/sbin/dnrd -k	
	eend $?
}

