# net-misc/zerotier/files/zerotier.init-r1
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="ZeroTier - Global Area Networking"
pidfile="/run/${SVCNAME}.pid"
command="/usr/sbin/zerotier-one"
command_background=true

depend() {
    use net
}
