# net-wireless/chillispot/files/chillispot
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
        use net logger
}


checkconfig() {
        if [ -f /etc/chilli.conf ]; then
                return 0;
        else
                eerror "Error starting chillispot. Please create /etc/chilli.conf before."
                return 1;
        fi
}

start() {
        checkconfig || return 1
        ebegin "Starting Chillispot"

        start-stop-daemon --start --quiet \
                --exec /usr/sbin/chilli
        eend $?
}

stop() {
        ebegin "Stopping Chillispot"
        start-stop-daemon --stop --pidfile /var/run/chilli.pid --quiet
        eend $?
}
