# net-dns/unbound/files/unbound.socket
[Socket]
ListenDatagram=127.0.0.1:1153
ListenStream=127.0.0.1:1153
[Install]
WantedBy=sockets.target
