# games-fps/rtcwmp-demo/files/rtcwmp-demo-ded.rc
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting RTCW demo dedicated server..."
	screen -A -m -d -S rtcwdemo su - dedicated -c GENTOO_DIR/startrtcwded
	eend $?
}

stop() {
	ebegin "Stopping RTCW demo dedicated server..."
	kill `screen -list | grep rtcwdemo | awk -F . '{ print $1 }' | sed -e s/.//`
	eend $?

}

status() {
	screen -list | grep rtcwdemoded
}
