# mail-filter/amavisd-new/files/amavis-mc.initd
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

pidfile="/run/amavis/amavis-mc.pid"
command="/usr/sbin/amavis-mc"
command_args="${amavis_mc_args} -P ${pidfile}"
start_stop_daemon_args="--interpreted"

depend() {
	use logger
	before amavisd-new snmpd
}
