# sys-power/acpid/files/systemd/acpid.socket
[Unit]
Description=ACPID Listen Socket

[Socket]
ListenStream=/var/run/acpid.socket

[Install]
WantedBy=sockets.target
