# net-proxy/cntlm/files/cntlm.service
[Unit]
Description=CNTLM HTTP Accelerator For NTLM Secured Proxies Authenticator
After=network.target

[Service]
ExecStart=/usr/bin/cntlm -f -U cntlm -c /etc/cntlm.conf

[Install]
WantedBy=multi-user.target
