# net-dns/dnssec-tools/files/donutsd.service
[Unit]
Description=DNSSEC monitoring daemon

[Service]
ExecStart=/usr/bin/donutsd -i /etc/dnssec-tools/donutsd.db

[Install]
WantedBy=multi-user.target
