# net-dns/rbldnsd/files/initd
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
	need net 
}

start() {
	ebegin "Starting rbldnsd-server"
	start-stop-daemon --start --quiet --pidfile /var/run/rbldnsd.pid --exec /usr/sbin/rbldnsd -- ${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping rbldnsd-server"
	start-stop-daemon --stop --quiet --pidfile /var/run/rbldnsd.pid --exec /usr/sbin/rbldnsd
	eend $?
}

