# net-proxy/shadowsocks-libev/files/shadowsocks-libev-redir_at.service
[Unit]
Description=Shadowsocks-Libev Client Service Redir Mode for %I
After=network.target

[Service]
Type=simple
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/ss-redir -c /etc/shadowsocks-libev/%i.json

[Install]
WantedBy=multi-user.target
