# app-admin/sshguard/files/sshguard.initd-r2
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/sbin/sshguard"
pidfile="${SSHGUARD_PIDFILE:-/var/run/${SVCNAME}.pid}"
command_args="-i \"${pidfile}\" ${SSHGUARD_OPTS}"
start_stop_daemon_args="--background --wait ${SSHGUARD_WAIT:-999} --interpreted --quiet"

depend() {
	after iptables
	use logger
}
