# sys-apps/systemd-tmpfiles/files/stmpfiles-setup.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Released under the 2-clause BSD license.

description="Set up tmpfiles.d entries"

depend()
{
	provide stmpfiles-setup stmpfiles.setup tmpfiles-setup tmpfiles.setup
	need localmount
}

start()
{
	ebegin "Setting up tmpfiles.d entries"
	/bin/systemd-tmpfiles --exclude-prefix=/dev --create --remove --boot \
		${stmpfiles_opts}
	eend $?
}
