# net-analyzer/openvas-scanner/files/gvm-feed-sync.cron
0 2 * * * gvm [ -x /etc/gvm/gvm-feed-sync.sh ] && /bin/bash /etc/gvm/gvm-feed-sync.sh > /dev/null
