# net-analyzer/openvas-scanner/files/openvassd.service
[Unit]
Description=Open Vulnerability Assessment Scanner
After=network.target
After=redis.service
Before=gvmd.service
Requires=redis.service

[Service]
Type=forking
EnvironmentFile=-/etc/openvas/sysconfig/openvassd-daemon.conf
ExecStart=/usr/bin/openvassd $OPENVAS_SCANNER_OPTIONS $OPENVAS_SCANNER_LISTEN_SOCKET $OPENVAS_SCANNER_LISTEN_OWNER $OPENVAS_SCANNER_LISTEN_GROUP $OPENVAS_SCANNER_LISTEN_MODE
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
