# net-p2p/gridcoin/files/gridcoin-staging.service
[Unit]
Description=Gridcoin Wallet Daemon (Staging)
After=network.target

[Service]
Type=forking
User=gridcoin
Group=gridcoin
ExecStart=/usr/bin/gridcoinresearchd-staging -testnet

[Install]
WantedBy=multi-user.target
