# net-proxy/c-icap/files/c-icap.init.2
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

pidfile=$(awk '$1 == "PidFile" { print $2 }' /etc/c-icap/c-icap.conf)

command="/usr/libexec/c-icap"
command_arguments="${EXTRA_OPTS}"

depend() {
    need localmount
}

start_pre() {
	cmdsocket=$(awk '$1 == "CommandsSocket" { print $2 }' /etc/c-icap/c-icap.conf)

	checkpath -d "$(dirname "${cmdsocket}")" "$(dirname "${pidfile}")"
}
