# net-misc/minidlna/files/minidlna-1.1.0.service
[Unit]
Description=MiniDLNA is a DLNA/UPnP-AV server software
After=syslog.target local-fs.target network.target

[Service]
User=minidlna
Group=minidlna
Type=forking
PIDFile=/var/run/minidlna/minidlna.pid
ExecStart=/usr/sbin/minidlnad -f /etc/minidlna.conf -P /var/run/minidlna/minidlna.pid

[Install]
WantedBy=multi-user.target
