# net-im/gg-transport/files/jggtrans-2.2.4
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	after net jabber-server
}

start() {
	ebegin "Starting Jabber GG Transport"
	start-stop-daemon --start --user jabber --group jabber --exec /usr/sbin/jggtrans -- -u jabber -g jabber /etc/jabber/jggtrans.xml
	eend $?
}

stop() {
	ebegin "Stopping Jabber GG Transport"
	start-stop-daemon --stop --pidfile /var/run/jabber/jggtrans.pid
	eend $?
}
