# dev-libs/cyrus-sasl/files/pwcheck.service
[Unit]
Description=SASL pwcheck daemon

[Service]
Type=forking
ExecStart=/usr/sbin/pwcheck

[Install]
WantedBy=multi-user.target
