# dev-libs/cyrus-sasl/files/pwcheck.rc6
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-libs/cyrus-sasl/files/pwcheck.rc6,v 1.6 2011/08/05 07:21:46 eras Exp $

depend() {
	need localmount
	use logger
}

start() {
	ebegin "Starting sasl pwcheck daemon"
	start-stop-daemon --start --quiet --exec /usr/sbin/pwcheck
	eend $?
}

stop() {
	ebegin "Stopping sasl pwcheck daemon"
	start-stop-daemon --stop --quiet --exec /usr/sbin/pwcheck
	eend $?
}
