# x11-base/xorg-server/files/xdm-setup.initd-1
#!/sbin/openrc-run
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need localmount
}

start() {
	if get_bootparam "nox" ; then
		touch /etc/.noxdm
	fi
}
