# net-firewall/nftables/files/systemd/nftables.service
[Unit]
Description=Store and restore nftables firewall rules

[Install]
Also=nftables-store.service
Also=nftables-restore.service
