# app-emulation/lxd/files/lxd-4.0.4-systemd-cgroups.initd
#!/sbin/openrc-run

depend() {
	need cgroups
}

start() {
	ebegin "Mounting systemd cgroups for openrc"
	mkdir -p /sys/fs/cgroup/systemd
	mount -t cgroup -o none,name=systemd systemd /sys/fs/cgroup/systemd
	eend $? "Failed to mount systemd cgroups"
}

stop() {
	ebegin "Unmounting systemd cgroups"
	umount -R /sys/fs/cgroup/systemd
	eend $? "Failed to unmount systemd cgroups"
}
