# app-emulation/lxd/files/lxd.service
[Unit]
Description=Container hypervisor based on LXC
After=cgmanager.service
Requires=cgmanager.service

[Service]
ExecStart=/usr/sbin/lxd --group lxd
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
