# net-analyzer/ndoutils/files/ndo2db.init-nagios3
#!/sbin/openrc-run
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depends() {
	before nagios
	need mysql
}

start() {
	ebegin "Starting ndo2db"
	if [ -S /var/nagios/ndo.sock ] ; then
		rm -f /var/nagios/ndo.sock
	fi
	start-stop-daemon --start --quiet --exec /usr/bin/ndo2db \
	-- -c /etc/nagios/ndo2db.cfg
	eend $?
}

stop() {
	ebegin "Stopping ndo2db"
	start-stop-daemon --stop --quiet --exec /usr/bin/ndo2db
	eend $?
}
