# media-sound/upmpdcli/files/upmpdcli.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="UPnP Media Renderer front-end for MPD, the Music Player Daemon"
pidfile=/var/run/upmpdcli.pid
command=/usr/bin/upmpdcli
config="${config:-/etc/upmpdcli.conf}"
command_args="-c ${config}"
command_args_background="-D"

depend() {
    need localmount mpd
}
