# sys-apps/pcsc-lite/files/pcscd-init.5
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

name="PC/SC Daemon"

pidfile=/var/run/pcscd/pcscd.pid

command=/usr/sbin/pcscd
command_args="${EXTRA_OPTS}"

start_stop_daemon_args="--user pcscd:pcscd"

depend() {
	need localmount
	after udev openct
	use logger
}

start_pre() {
	checkpath -d -m 0755 -o pcscd:pcscd /var/run/pcscd
}
