# app-misc/linux-logo/files/linux-logo.service
[Unit]
Description=linux-logo-creator
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/linux_logo_creator

[Install]
WantedBy=logo.target
