# sci-misc/boinc/files/boinc.service
[Unit]
Description=BOINC Daemon

[Service]
User=boinc
Group=boinc
Nice=19
ExecStart=/usr/bin/boinc_client --dir /var/lib/boinc --redirectio

[Install]
WantedBy=multi-user.target
