# net-misc/olsrd/files/olsrd-r1.initd
#!/sbin/openrc-run
command="/usr/sbin/olsrd"
command_background="yes"
pidfile="/run/olsrd.pid"
