# dev-perl/MogileFS-Server/files/mogautomount-init.d-2.50
#!/sbin/openrc-run

NAME="mogautomount"
BINARY="/usr/bin/mogautomount"

depend() {
	use net
	before mogstored
}

start() {
	ebegin "Mounting MogileFS disks"
	${BINARY} --chmod-mountpoints
	eend $?
}

stop() {
	:
}

# vim: ft=gentoo-init-d:
