# dev-perl/MogileFS-Server/files/mogilefsd-init.d-2.66
#!/sbin/openrc-run
# $Id$

NAME="mogilefsd"
BINARY="/usr/bin/mogilefsd"
MOGILE_UID="mogile"
MOGILE_GID="mogile"

depend() {
	use net mysql postgresql
	after mysql postgresql
}

start() {
	install -o $MOGILE_UID -g $MOGILE_GID -d $(dirname $PIDFILE)
	ebegin "Starting $NAME"
	start-stop-daemon --chuid ${MOGILE_UID}:${MOGILE_GID} --start \
		--pidfile ${PIDFILE} --exec ${BINARY} \
		-- ${MOGILEFSD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping $NAME"
	start-stop-daemon --stop --quiet --pidfile="${PIDFILE}"
	eend $?
}

# vim: ft=gentoo-init-d:
