# app-laptop/samsung-tools/files/samsung-tools.init
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need dbus
}

start() {
	ebegin "Setting up samsung-tools"
	dbus-send --system --print-reply=literal --dest='org.voria.SamsungTools.System' \
		'/' org.voria.SamsungTools.System.SetInitialDevicesStatus
	eend $?
}
