# net-misc/rsync/files/rsyncd.service
[Unit]
Description=rsync daemon
After=networking.target

[Service]
Type=simple
ExecStart=/usr/bin/rsync --daemon --no-detach
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
