# net-nds/openldap/files/slapd-2.6.1.service
[Unit]
Description=OpenLDAP Server Daemon
After=network.target

[Service]
Type=notify
PIDFile=/run/openldap/slapd.pid
ExecStartPre=/usr/sbin/slaptest -Q -u $SLAPD_OPTIONS
ExecStart=/usr/lib/openldap/slapd -u ldap -h ${SLAPD_URLS} $SLAPD_OPTIONS
NotifyAccess=all

[Install]
WantedBy=multi-user.target
