# net-nds/openldap/files/slurpd-initd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-nds/openldap/files/slurpd-initd,v 1.1 2007/01/16 23:22:02 jokey Exp $

depend() {
	need net
}

start() {
	ebegin "Starting slurpd"
	start-stop-daemon --start --quiet \
	--exec /usr/lib/openldap/slurpd 
	eend $?
}

stop() {
	ebegin "Stopping slurpd"
	start-stop-daemon --stop --quiet \
	--exec /usr/lib/openldap/slurpd 
	eend $?
}
