# sys-fs/zfs-auto-snapshot/files/zfs-auto-snapshot-weekly.timer
[Unit]
Description=ZFS weekly snapshot timer
Documentation=man:zfs-auto-snapshot(8)

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=600
