# sys-fs/zfs-auto-snapshot/files/zfs-auto-snapshot-daily.timer
[Unit]
Description=ZFS daily snapshot timer
Documentation=man:zfs-auto-snapshot(8)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600
