# sys-fs/zfs-auto-snapshot/files/zfs-auto-snapshot-hourly.timer
[Unit]
Description=ZFS hourly snapshot timer
Documentation=man:zfs-auto-snapshot(8)

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=59
