# sys-fs/zfs-auto-snapshot/files/zfs-auto-snapshot-frequent.service
[Unit]
Description=ZFS frequent snapshot service
Documentation=man:zfs-auto-snapshot(8)

[Service]
ExecStart=@GENTOO_PORTAGE_EPREFIX@/usr/sbin/zfs-auto-snapshot @DEFAULT_EXCLUDE@ --skip-scrub --label=frequent --keep=4 //
