# net-misc/sntpd/files/sntpd.initd-r1
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="sntpd"
pidfile="/run/sntpd.pid"

command="/usr/sbin/sntpd"
command_args="-n ${SNTPD_OPTS}"
command_background="true"

depend() {
	use net dns logger
}
