# net-irc/inspircd/files/inspircd-init
#!/sbin/openrc-run
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

extra_commands="rehash version"

depend() {
	need net
	provide ircd
}

start() {
	ebegin "Starting InspIRCd"
	checkpath -d -m 0750 -o inspircd:inspircd /var/run/inspircd
	start-stop-daemon --start --quiet --user inspircd \
		--exec /usr/bin/inspircd -- \
		--config /etc/inspircd/inspircd.conf \
		--logfile /var/log/inspircd/ircd.log
	eend $?
}

stop() {
	ebegin "Stopping InspIRCd"
	start-stop-daemon --stop --quiet --exec /usr/bin/inspircd
	eend $?
}

rehash() {
	ebegin "Rehashing InspIRCd"
	/usr/lib/inspircd/inspircd.launcher/inspircd rehash
	eend $?
}

version() {
	ebegin "Retrieve InspIRCd version"
	/usr/lib/inspircd/inspircd.launcher/inspircd version
	eend $?
}
