# sys-apps/man-db/files/man-db.timer
[Unit]
Description=Daily man-db cache update

[Timer]
OnCalendar=daily
AccuracySec=1d
Persistent=true

[Install]
WantedBy=timers.target
