# sys-apps/man-db/files/man-db.cron
#!/bin/sh

# Use same perms/settings as the ebuild.
if [ ! -d /var/cache/man ]; then
	mkdir -p /var/cache/man
	chown man:root /var/cache/man
	chmod 2755 /var/cache/man
fi

exec nice mandb --quiet
