# www-servers/hinsightd/files/confd-v1.sh
#!/sbin/openrc-run

#HIN_NAME=hinsightd
#HIN_RUN_DIR=/var/run/$HIN_NAME
#HIN_LOG_DIR=/var/log/$HIN_NAME
#HIN_TMP_DIR=/var/tmp/$HIN_NAME
#HIN_WORK_DIR=/var/www/localhost
#HIN_CFG_FILE=/etc/hinsightd/main.lua
#HIN_LOG_FILE=/var/log/hindsightd/hindsightd.log
#HIN_BIN=/usr/sbin/hinsightd
#HIN_USER=hinsightd

