# www-servers/hinsightd/files/hinsightd.service
[Unit]
Description=hinsightd is a http/1.1 webserver
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/hinsightd --config /etc/hinsightd/hinsightd.lua
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target
