# net-misc/sslh/files/sslh.init.d-3
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Port multiplexer for ssh/ssl/etc..."
command="sslh"
pidfile="/run/${SVCNAME}.pid"
command_args="${DAEMON_OPTS} --pidfile ${pidfile}"
