# sys-power/apcupsd/files/apcupsd.service
[Unit]
Description=APC UPS Monitor

[Service]
ExecStart=/sbin/apcupsd -b

[Install]
WantedBy=multi-user.target
