# sys-power/apcupsd/files/apcupsd.service
[Unit]
Description=APC UPS Monitor
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/sbin/apcupsd -b

[Install]
WantedBy=multi-user.target
