# app-misc/beanstalkd/files/init-1.4.6
#!/sbin/openrc-run
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting beanstalkd"
	/sbin/start-stop-daemon --start \
        --background \
		--pidfile ${PIDFILE} --make-pidfile \
		--exec ${BEANSTALKD_BINARY} \
		-- -b ${DATADIR} -p ${PORT} -l ${ADDR} -u ${USER}
	eend $?
}

stop() {
	ebegin "Stopping beanstalkd"
	start-stop-daemon --stop --quiet \
		--pidfile ${PIDFILE} \
		--exec ${BEANSTALKD_BINARY}
	eend $?
}
