# net-vpn/miredo/files/miredo.rc.2
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need localmount
	after net modules
	use dns logger
}

command=/usr/sbin/${SVCNAME}
command_args="${MIREDO_OPTS}"
pidfile=/var/run/${SVCNAME}.pid
