# app-misc/lirc/files/lircmd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-misc/lirc/files/lircmd,v 1.2 2004/09/28 00:22:15 swegener Exp $

depend() {
	need lircd
}

start() {
	ebegin "Starting lircmd"
	start-stop-daemon --start --quiet --exec /usr/sbin/lircmd
	eend $?
}

stop() {
	ebegin "Stopping lircmd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/lircmd
	eend $?
}
