# net-p2p/deluge/files/deluge-web.service
[Unit]
Description=Deluge WebUI
After=deluged.service

[Service]
User=deluge
ExecStart=/usr/bin/deluge-web

[Install]
WantedBy=multi-user.target

