# app-backup/boxbackup/files/bbackupd.rc
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need localmount
	use net
}

start() {
	ebegin "Starting box backup daemon"
	start-stop-daemon --start --exec /usr/sbin/bbackupd > /dev/null
	eend $?
}
 
stop() {
	ebegin "Stopping box backup daemon"
	start-stop-daemon --stop --exec /usr/sbin/bbackupd
	eend $?
}
