# net-mail/courier-imap/files/courier-authdaemond-r1.service
[Unit]
Description=Courier Authentification Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/courierlogger -pid=/var/run/authdaemon.pid -start /usr/libexec/courier-authlib/authdaemond
ExecStop=/usr/sbin/courierlogger -pid=/var/run/authdaemon.pid -stop /usr/libexec/courier-authlib/authdaemond
PIDFile=/var/run/authdaemon.pid

[Install]
WantedBy=multi-user.target
