# net-analyzer/openvas-manager/files/openvasmd.service
[Unit]
Description=OpenVAS Manager
After=network.target
After=openvassd.service

[Service]
EnvironmentFile=-/etc/openvas/openvasmd.conf
ExecStart=/usr/sbin/openvasmd -f ${MANAGER_LISTEN} ${MANAGER_PORT} ${SCANNER_LISTEN} ${SCANNER_PORT} ${MANAGER_EXTRA_ARGS}
Restart=always
RestartSec=1
User=root
Group=root
TimeoutSec=1200

[Install]
WantedBy=multi-user.target
