# net-misc/wide-dhcpv6/files/dhcp6c-AT.service
[Unit]
Description=Wide DHCPv6 client daemon
Documentation=man:dhcp6c(8)
After=sys-subsystem-net-devices-%i.device
Requires=sys-subsystem-net-devices-%i.device

[Service]
ExecStart=dhcp6c -f -p /run/dhcp6c-%i.pid %I
ExecReload=/bin/kill -HUP $MAINPID
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
