# net-misc/wide-dhcpv6/files/wide-dhcp6s.init
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="Wide DHCPv6 server"
command="/usr/sbin/dhcp6s"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="${OPTIONS} ${INTERFACE}"
command_args_background="-P ${pidfile}"

depend() {
	need net
}
