# media-sound/murmur/files/murmurd-no-dbus.service
[Unit]
Description=Murmurd - VoIP Server (for Mumble)
After=network.target mysql.service

[Service]
User=murmur
Group=murmur
ExecStart=/usr/bin/murmurd -fg -ini /etc/murmur/murmur.ini

[Install]
WantedBy=multi-user.target
