# app-crypt/ekeyd/files/ekeyd.service
[Unit]
Description=Entropy key daemon

[Service]
ExecStart=/usr/libexec/ekeyd 

[Install]
WantedBy=multi-user.target
