# net-analyzer/snmptt/files/snmptt.initd-r1
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need snmptrapd
}

start() {
	ebegin "Starting snmptt"
	start-stop-daemon --start --pidfile /run/snmptt.pid --exec /usr/sbin/snmptt
	eend
}

stop() {
	ebegin "Stopping snmptt"
	start-stop-daemon --stop --pidfile /run/snmptt.pid
	eend
}
