# app-shells/atuin/files/atuin.service
[Unit]
Description=Atuin (shell history manager) server
After=network.target

[Service]
ExecStart=/usr/bin/atuin server start
User=atuin

[Install]
WantedBy=default.target
