# app-admin/logcheck/files/logcheck.timer
[Unit]
Description=Mails anomalies in the system logfiles to the administrator

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
