# app-admin/logcheck/files/logcheck.cron.d
@reboot sleep 60; /etc/cron.hourly/logcheck -R
