# media-video/motion/files/motion.service
[Unit]
Description=Motion detection daemon
After=local-fs.target network.target
Documentation=man:motion(1)

[Service]
Type=simple
ExecStart=/usr/bin/motion -n -c /etc/motion/motion.conf
StandardError=null

#User=motion
Group=video

#RuntimeDirectory=motion
#RuntimeDirectoryMode=0750

#WorkingDirectory=/var/lib/motion

[Install]
WantedBy=multi-user.target
