# net-proxy/tinyproxy/files/tinyproxy-1.10.0.service
[Unit]
Description=Tinyproxy Web Proxy Server
After=network.target

[Service]
ExecStart=/usr/bin/tinyproxy -d -c /etc/tinyproxy/tinyproxy.conf

[Install]
WantedBy=multi-user.target
