# net-proxy/tinyproxy/files/tinyproxy.service
[Unit]
Description=Tinyproxy Web Proxy Server
After=network.target

[Service]
ExecStart=/usr/sbin/tinyproxy -d -c /etc/tinyproxy.conf

[Install]
WantedBy=multi-user.target
