# media-radio/ax25-apps/files/ax25ipd.rc
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/media-radio/ax25-apps/files/ax25ipd.rc,v 1.3 2004/07/14 22:42:51 agriffis Exp $

depend() {
	need net
}

start() {
	ebegin "Starting ax25ipd"
	start-stop-daemon --start --quiet --background --exec /usr/sbin/ax25ipd
	eend $?
}

stop() {
	ebegin "Stopping ax25ipd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/ax25ipd
	eend $?
}
