# net-wireless/unifi/files/unifi.service
[Unit]
Description=UniFi Controller
Requires=network.target
After=network.target

[Service]
ExecStart=/usr/bin/java -jar /usr/lib/unifi/lib/ace.jar start
ExecStop=/usr/bin/java -jar /usr/lib/unifi/lib/ace.jar stop
User=unifi
Group=unifi
WorkingDirectory=/usr/lib/unifi

[Install]
WantedBy=network.target
