# www-apps/grafana-bin/files/grafana.initd.3
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need localmount
}

start() {
	ebegin "Starting grafana server"

    checkpath -d -o "${GRAFANA_USER}:${GRAFANA_GROUP}" -m750 "${GRAFANA_DATADIR}"
    checkpath -d -o "${GRAFANA_USER}:${GRAFANA_GROUP}" -m750 "${GRAFANA_DATADIR}/dashboards"
    checkpath -d -o "${GRAFANA_USER}:${GRAFANA_GROUP}" -m750 "${GRAFANA_DATADIR}/plugins"
    checkpath -d -o "${GRAFANA_USER}:${GRAFANA_GROUP}" -m750 "${GRAFANA_DATADIR}/sessions"
    checkpath -d -o "${GRAFANA_USER}:${GRAFANA_GROUP}" -m750 "${GRAFANA_LOGDIR}"

	start-stop-daemon --start -b \
		--make-pidfile \
		--pidfile /run/${SVCNAME}.pid \
		--chdir /usr/share/grafana \
		--user ${GRAFANA_USER}:${GRAFANA_GROUP} \
		--exec /usr/bin/grafana-server -- \
		-config=/etc/grafana/${SVCNAME}.ini \
		cfg:default.paths.data=${GRAFANA_DATADIR} \
		cfg:default.paths.logs=${GRAFANA_LOGDIR}

	eend $?
}

stop() {
	ebegin "Stopping grafana server"

	start-stop-daemon --stop \
		--pidfile /run/${SVCNAME}.pid \
		--user "${GRAFANA_USER}" \
		--retry=TERM/20/KILL/5

	eend $?
}
