# www-apps/grafana-bin/files/grafana.initd.3
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

user="${GRAFANA_USER:-grafana}"
group="${GRAFANA_GROUP:-grafana}"
datadir="${GRAFANA_DATADIR:-/var/lib/grafana}"
logdir="${GRAFANA_LOGDIR:-/var/log/grafana}"

depend() {
	need localmount
}

start() {
	ebegin "Starting grafana server"

    checkpath -d -o "${user}:${group}" -m750 "${datadir}"
    checkpath -d -o "${user}:${group}" -m750 "${datadir}/dashboards"
    checkpath -d -o "${user}:${group}" -m750 "${datadir}/plugins"
    checkpath -d -o "${user}:${group}" -m750 "${datadir}/sessions"
    checkpath -d -o "${user}:${group}" -m750 "${logdir}"

	start-stop-daemon --start -b \
		--make-pidfile \
		--pidfile /run/${SVCNAME}.pid \
		--chdir /usr/share/grafana \
		--user ${user}:${group} \
		--exec /usr/bin/grafana-server -- \
		-config=/etc/grafana/${SVCNAME}.ini \
		cfg:default.paths.data=${datadir} \
		cfg:default.paths.logs=${logdir}

	eend $?
}

stop() {
	ebegin "Stopping grafana server"

	start-stop-daemon --stop \
		--pidfile /run/${SVCNAME}.pid \
		--user "${user}" \
		--retry=TERM/20/KILL/5

	eend $?
}
