# net-voip/yate/files/yate.rc6
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
    need net
}

start() {
    ebegin "Starting YATE"
    start-stop-daemon --start --exec /usr/bin/yate -- -d -p /var/run/yate.pid ${YATE_OPTS}
    eend $?
}

stop() {
    ebegin "Stopping YATE"
    start-stop-daemon --stop --pidfile=/var/run/yate.pid
    eend $?
}
