# sys-process/atop/files/atop.rc-r1
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-process/atop/files/atop.rc-r1,v 1.1 2015/06/25 21:13:42 polynomial-c Exp $

description="Resource-specific view of processes"
pidfile="/var/run/atop.pid"
command="/usr/bin/atop"

start() {
	ebegin "Starting atop"
	start-stop-daemon --start --pidfile ${pidfile} \
		--exec /etc/atop/atop.daily
	eend $?
}
