# sys-process/atop/files/atop.rc
#!/sbin/openrc-run
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

start() {
	ebegin "Starting atop"
	start-stop-daemon --start --quiet --exec /etc/atop/atop.daily
	eend $?
}

stop() {
	ebegin "Stopping atop"
	start-stop-daemon --stop --exec /usr/bin/atop --pidfile /var/run/atop.pid
	eend $?
}
