# app-crypt/mit-krb5/files/mit-krb5kpropd.service
[Unit]
Description=Kerberos 5 propagation server

[Service]
ExecStart=/usr/sbin/kpropd -S

[Install]
WantedBy=multi-user.target
