# app-crypt/mit-krb5/files/mit-krb5kpropd_at.service
[Unit]
Description=Kerberos 5 propagation server
Conflicts=mit-krb5kpropd.service

[Service]
ExecStart=/usr/sbin/kpropd
StandardInput=socket
