# app-crypt/mit-krb5/files/mit-krb5kpropd.socket
[Unit]
Description=Kerberos 5 propagation server

[Socket]
ListenStream=754
Accept=yes

[Install]
WantedBy=sockets.target
