# dev-libs/openct/files/openct.rc
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need localmount
}

checkconfig() {
	if [ ! -e /etc/openct.conf ] ; then
		eerror "You need an /etc/openct.conf file to run OpenCT"
		return 1
	fi
	return 0
}

start() {
	checkconfig || return 1

	ebegin "Starting OpenCT"
	/usr/sbin/openct-control init
	eend $?
}

stop() {
	ebegin "Stopping OpenCT"
	/usr/sbin/openct-control shutdown
	eend $?
}
