# app-admin/rasdaemon/files/rasdaemon.openrc
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Starts ${SVCNAME} service"

command="/usr/sbin/${SVCNAME}"
command_args="--foreground --record"
pidfile=/run/${SVCNAME}.pid
command_background=true

start_post() {
  "${command}" --enable
}

stop_post() {
  "${command}" --disable
}
