# sys-process/below/files/below.service
[Unit]
Description=below system monitor recording daemon
After=time-sync.target

[Service]
ExecStart=/usr/bin/below record --retain-for-s 604800 --compress
# Enable backtraces in errors
Environment=RUST_LIB_BACKTRACE=1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
