# net-fs/samba/files/nmbd.service
[Unit]
Description=Samba NetBIOS name server
After=network.target

[Service]
Type=forking
PIDFile=/var/run/samba/nmbd.pid
ExecStart=/usr/sbin/nmbd -D 
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
