# net-fs/samba/files/smbd.service
[Unit]
Description=Samba SMB/CIFS server
After=network.target nmbd.service winbindd.service

[Service]
Type=forking
PIDFile=/var/run/samba/smbd.pid
ExecStart=/usr/sbin/smbd -D 
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
