# sys-process/vixie-cron/files/vixie-cron.rc7
#!/sbin/openrc-run
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command=/usr/sbin/cron
pidfile=/var/run/cron.pid

depend() {
	use clock logger
	need localmount
	provide cron
}
