# sys-process/vixie-cron/files/vixie-cron.service
[Unit]
Description=Vixie Cron Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/cron
ExecStop=/bin/kill -TERM $MAINPID
IgnoreSIGPIPE=false

[Install]
WantedBy=multi-user.target
