# sys-process/vixie-cron/files/vixie-cron.rc7
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-process/vixie-cron/files/vixie-cron.rc7,v 1.1 2011/12/31 14:58:21 idl0r Exp $

command=/usr/sbin/cron
pidfile=/var/run/cron.pid

depend() {
	use clock logger
	need localmount
	provide cron
}
