# sys-apps/audio-entropyd/files/audio-entropyd.service
[Unit]
Description=Audio Entropy Daemon

[Service]
ExecStart=/usr/sbin/audio-entropyd -n

[Install]
WantedBy=multi-user.target
