# app-admin/fluent-bit/files/initd-1
#!/sbin/openrc-run
# Copyright 2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="fluent-bit daemon"
description=""
command=/usr/bin/fluent-bit
command_args="-c ${FLUENT_CONF:-/etc/fluent-bit/fluent-bit.conf}"
command_user="fluentd"
command_group="fluentd"
command_background=yes
pidfile="/run/fluent-bit.pid"

depend() {
	need net localmount
}
