# sys-power/cpupower/files/cpupower-frequency-set.service
[Unit]
Description=CPU Frequency Scaling

[Service]
Type=oneshot
ExecStart=/usr/bin/cpupower frequency-set $OPTS

[Install]
WantedBy=multi-user.target
