# x11-misc/afbinit/files/afbinit.init
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	before xdm
}

start() {
	ebegin "Loading afbinit"
	/usr/sbin/rc.afb
	eend $? "Failed to load afbinit"
}
