# net-vpn/cloudflare-warp/files/warp-svc.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Cloudflare Zero Trust Client Daemon"
supervisor="supervise-daemon"
command="/bin/warp-svc"
command_args=""

depend() {
	need net
}
