# net-misc/openntpd/files/openntpd.init.d-20080406-r6
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

name="OpenNTPD"
command="/usr/sbin/ntpd"
pidfile="/run/ntpd.pid"
command_args="-p ${pidfile} ${NTPD_OPTS}"

depend() {
	need net
	after ntp-client
	use dns logger
}
