# app-backup/backuppc/files/backuppc.service
[Unit]
Description= BackupPC server
After=syslog.target local-fs.target remote-fs.target

[Service]
Type=oneshot
User=backuppc
Group=backuppc
ExecStart=/usr/bin/BackupPC -d
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
