# sys-apps/mlocate/files/updatedb.service
[Unit]
Description=Update locate database

[Service]
Type=oneshot
ExecStart=/usr/bin/updatedb
IOSchedulingClass=idle
