# mail-filter/policyd/files/policyd-cleanup.cron
#!/bin/bash
/usr/sbin/policyd_cleanup -c /etc/policyd.conf
