# sys-freebsd/freebsd-cddl/files/zfs.initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need hostid
	before fsck
}

start() {
	einfo "Starting zfs ..."
	zfs mount -a
	zfs share -a
	if [ ! -r /etc/zfs/exports ]; then
		touch /etc/zfs/exports
	fi
}

stop() {
	einfo "Stopping zfs ..."
	zfs unshare -a
	zfs unmount -a
}
