# net-nntp/inn/files/innd
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting innd"
	checkpath -d -m 0755 -o news:news /var/run/news
	start-stop-daemon --start --user news --exec /usr/lib/news/bin/rc.news
	eend $?
}

stop() {
	ebegin "Stopping innd"
	start-stop-daemon --start --user news --exec /usr/lib/news/bin/rc.news -- stop
	sleep 2
	eend $?
}
