# net-misc/sobby/files/sobby-init-0.4.7
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2         
# $Header: /var/cvsroot/gentoo-x86/net-misc/sobby/files/sobby-init-0.4.7,v 1.2 2011/12/10 09:30:34 xarthisius Exp $

depend() {
	use avahi
}

start() {
	ebegin "Starting obby dedicated server"
	start-stop-daemon -b --quiet --start  --make-pidfile \
		--pidfile /var/run/sobby.pid -u ${SOBBY_USER} -g ${SOBBY_GROUP} \
		--exec /usr/bin/sobby -- ${SOBBY_OPTS}
	eend $? "Failed to start sobby"
}

stop() {
	ebegin "Stopping obby dedicated server"
	start-stop-daemon --stop --quiet --pidfile /var/run/sobby.pid
	eend $? "Failed to stop sobby"
}
