# net-dialup/rp-l2tp/files/rp-l2tpd-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting rp-l2tpd"
	start-stop-daemon --start --quiet --exec /usr/sbin/rp-l2tpd
	eend $?
}

stop() {
	ebegin "Stopping rp-l2tpd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/rp-l2tpd
	eend $?
}
