# sys-process/cronie/files/cronie-1.3-initd
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-process/cronie/files/cronie-1.3-initd,v 1.3 2014/03/04 00:17:58 polynomial-c Exp $

command="/usr/sbin/crond"
command_args="${CRONDARGS}"
pidfile="/var/run/crond.pid"

depend() {
	use clock logger
	need localmount
	provide cron	
}
