# net-voip/openmcu/files/openmcu.confd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

OPTIONS="-u openmcu -g openmcu"
