# sys-power/powernowd/files/powernowd.service
[Unit]
Description=PowerNow! Daemon

[Service]
ExecStart=/usr/sbin/powernowd -d -q

[Install]
WantedBy=multi-user.target
