# net-analyzer/pnp4nagios/files/npcd.initd
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

extra_started_commands="reload"

command="/usr/bin/npcd"
command_args="-f /etc/pnp/npcd.cfg -d"

depend() {
	need net
}

reload() {
        ebegin "Reloading npcd"
        kill -HUP `pgrep npcd`
        eend $? "Failed to reload npcd"
}
