# net-firewall/lutelwall/files/lutelwall
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/net-firewall/lutelwall/files/lutelwall,v 1.1 2005/03/27 18:11:00 vanquirius Exp $

depend() {
	need net
	use logger
}

start() {
	ebegin "Starting LutelWall"
	/usr/sbin/lutelwall start
	eend $?
}

stop() {
	ebegin "Stopping LutelWall"
	/usr/sbin/lutelwall stop
	eend $?
}

restart() {
	stop
	start
}
