# net-firewall/lutelwall/files/lutelwall
#!/sbin/openrc-run
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Id$

depend() {
	need net
	use logger
}

start() {
	ebegin "Starting LutelWall"
	/usr/sbin/lutelwall start
	eend $?
}

stop() {
	ebegin "Stopping LutelWall"
	/usr/sbin/lutelwall stop
	eend $?
}

restart() {
	stop
	start
}
