# www-servers/hiawatha/files/hiawatha.service
[Unit]
Description=Hiawatha webserver
After=network.target

[Service]
ExecStart=/usr/sbin/hiawatha -d

[Install]
WantedBy=multi-user.target
