# dev-db/hyperdex/files/hyperdex.initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-db/hyperdex/files/hyperdex.initd,v 1.1 2012/08/18 14:22:54 patrick Exp $

command="/usr/bin/hyperdex-daemon"
command_args="--daemon ${HYPERDEX_OPTS}"

depend() {
	use net
}
