# sys-process/at/files/atd.service
[Unit]
Description=at (job queue) daemon
After=networking.target

[Service]
Type=forking
ExecStart=/usr/sbin/atd
PIDFile=/run/atd.pid
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
