# sys-process/at/files/atd.rc7
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

pidfile="/var/run/atd.pid"
command="/usr/sbin/atd"
command_args="${ATD_OPTS}"
start_stop_daemon_args="--pidfile ${pidfile}"

depend() {
	need clock logger
}
