# sys-process/at/files/atd.confd
#!/sbin/openrc-run
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

# Define startup options for atd(8).
# For example: ATD_OPTS="-l 3.0 -b 5"
ATD_OPTS=""
