# net-fs/nfs-utils/files/proc-fs-nfsd.mount
[Unit]
Description=NFS Configuration File System
DefaultDependencies=no

[Mount]
What=sunrpc
Where=/proc/fs/nfsd
Type=nfsd
