# net-fs/nfs-utils/files/rpc-svcgssd.service
[Unit]
Description=RPC GSS-API Server Daemon
Documentation=man:rpc.svcgssd(8)
After=rpcbind.service
Requires=rpcbind.service
Before=nfs-server.target
PartOf=nfs-server.target

[Service]
ExecStart=/usr/sbin/rpc.svcgssd -f

[Install]
WantedBy=multi-user.target
