# net-fs/nfs-utils/files/var-lib-nfs-rpc_pipefs.mount
[Unit]
Description=RPC Pipe File System

[Mount]
What=rpc_pipefs
Where=/var/lib/nfs/rpc_pipefs
Type=rpc_pipefs
