# games-fps/ut2003/files/ut2003-ded.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile="/var/run/${RC_SVCNAME}.pid"
command_background="true"
command_user="${ut2003_ded_user}:${ut2003_ded_group}"
command="ut2003-ded"
command_args="${ut2003_ded_opts}"

depend() {
	need net
}
