# net-print/lprng/files/lprng-init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

checkconfig() {
	if [ ! -e /etc/lprng/lpd.conf ] ; then
		eerror "You need an /etc/lprng/lpd.conf file first"
		eerror "There is a sample in  /usr/share/doc/LPRng"
		return 1
	fi
	if [ ! -e /etc/lprng/printcap ] ; then
		eerror "You need an /etc/lprng/printcap file first"
		eerror "There is a sample in  /usr/share/doc/LPRng"
		return 1
	fi
	if [ ! -e /etc/lprng/lpd.perms ] ; then
		eerror "You need an /etc/lprng/lpd.perms file first"
		eerror "There is  a sample in  /usr/share/doc/LPRng"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting lpd"
	start-stop-daemon --start --pidfile "/var/run/lpd.515" --quiet --exec /usr/sbin/lpd
	eend $?
}

stop() {
	ebegin "Stopping lpd"
	start-stop-daemon --stop --pidfile "/var/run/lpd.515" --quiet --exec /usr/sbin/lpd
	eend $?
}
