# app-admin/mcelog/files/mcelog.init-r1
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/mcelog/files/mcelog.init-r1,v 1.1 2014/09/09 01:58:47 hasufell Exp $

description='Start/stop mcelog in daemon mode'

pidfile="/var/run/mcelog.pid"
command="/usr/sbin/mcelog"
command_args="--daemon --pidfile ${pidfile} ${MCELOG_OPTS}"

depend() {
	after udev 
	need localmount
	use logging
}
