# app-admin/mcelog/files/mcelog.init-r1
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description='Start/stop mcelog in daemon mode'

pidfile="/var/run/mcelog.pid"
command="/usr/sbin/mcelog"
command_args="--daemon --pidfile ${pidfile} ${MCELOG_OPTS}"

depend() {
	after udev 
	need localmount
	use logging
}
