# app-antivirus/clamav/files/clamd.service
[Unit]
Description=clamd scanner daemon
After=nss-lookup.target network.target

[Service]
Type=forking
ExecStart=/usr/sbin/clamd
Restart=on-failure
PrivateTmp=true

[Install]
WantedBy=multi-user.target
