# app-antivirus/clamav/files/clamd.service
.include /usr/lib/systemd/system/clamd@.service

[Unit]
Description=Generic ClamAV scanner daemon

[Install]
WantedBy=multi-user.target
