# app-antivirus/clamav/files/freshclam.initd
#!/sbin/openrc-run

command="/usr/bin/freshclam"
pidfile="/run/${RC_SVCNAME}.pid"

# Ignore the value of "PidFile" set in freshclam.conf.
command_args="-p ${pidfile}"
command_args_background="--daemon"
