# www-apps/readarr/files/readarr.init
#!/sbin/openrc-run

pidfile="/run/readarr.pid"
command="/opt/readarr/Readarr"
command_user="readarr:readarr"
command_background=true
