# www-servers/monkeyd/files/monkeyd.initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/www-servers/monkeyd/files/monkeyd.initd,v 1.1 2012/06/07 08:22:06 blueness Exp $

DAEMON="/usr/bin/monkey"
CONFFILE="/etc/monkeyd/monkey.conf"

depend() {
	use net
}

start() {
	ebegin "Starting monkeyd"
	start-stop-daemon --start --exec "${DAEMON}" -- "-D" "${MONKEY_ARGS}" >/dev/null
	eend $?
}

stop() {
	ebegin "Stopping monkeyd"
	local PORT=$(awk '/^ *Port/ { print $2 }' "${CONFFILE}")
	local PIDFILE=$(awk '/^ *PidFile/ { print $2 }' "${CONFFILE}")
	PIDFILE="${PIDFILE}"."${PORT}"
	start-stop-daemon --stop --quiet --pidfile "${PIDFILE}"
	eend $?
}
