# sys-apps/util-linux/files/uuidd.initd
#!/sbin/openrc-run
# Copyright 2024 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command=/usr/sbin/uuidd
command_args_foreground="-FP"
command_user=uuidd
pidfile=/run/uuidd/uuidd.pid

depend() {
	need clock localmount
}

start_pre() {
	checkpath -d -m 2755 -o uuidd:uuidd /run/uuidd &&
	checkpath -d -m 0755 -o uuidd:uuidd /var/lib/libuuid
}
