# sci-mathematics/rstudio/files/rstudio-server.initd
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
	ebegin "Starting RStudio Server"
	start-stop-daemon --start --quiet --exec /usr/bin/rserver
	eend $?
}

stop() {
	ebegin "Stopping RStudio Server"
	start-stop-daemon --stop --quiet --retry=TERM/30/KILL/5 --name /usr/bin/rserver
	eend $?
}
