# media-sound/shoutcast-server-bin/files/shoutcast.2
#!/sbin/runscript

depend() {
	need localmount
}

checkconfig() {
	if ! [ -f /etc/shoutcast/sc_serv.conf ]; then
		eerror "Can't find /etc/shoutcast/sc_serv.conf!"
		return 1
	fi

	return 0
}

start() {
	checkconfig || return 1

        ebegin "Starting SHOUTcast"
	start-stop-daemon --background --start --quiet --user shoutcast \
		--group nogroup --make-pidfile --pidfile /var/run/sc_serv.pid \
		--exec /opt/bin/sc_serv /etc/shoutcast/sc_serv.conf
	eend $?
}

stop() {
	ebegin "Stopping SHOUTcast"
	start-stop-daemon --stop --pidfile /var/run/sc_serv.pid --name sc_serv
	eend $?
}
