# net-misc/rinetd/files/rinetd.rc
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	use net
}

start() {
	ebegin "Starting rinetd"
	/usr/sbin/rinetd
	eend $?
}

stop() {
	local ret
	ebegin "Stopping rinetd"
	start-stop-daemon --stop --pidfile /var/run/rinetd.pid
	ret=$?
	rm -f /var/run/rinetd.pid
	eend ${ret}
}
