# sys-apps/preload/files/preload-0.6.4.init.in
#!/sbin/openrc-run

PIDFILE="/var/run/@PACKAGE@.pid"

depend() {
	after localmount
	use netmount ntpd
}

dump() {
	ebegin "Dumping config and state for @PACKAGE@"
	kill -USR1 $(<${PIDFILE})
	kill -USR2 $(<${PIDFILE})
	eend $?
}

reload() {
	ebegin "Reloading @PACKAGE@"
	kill -HUP $(<${PIDFILE})
	eend $?
}

start() {
	ebegin "Starting @PACKAGE@"

	#MIN_MEMORY=${MIN_MEMORY:-256}
	# Check for > MIN_MEMORY MB
	#free -m | awk '/Mem:/ {exit ($2 >= ('"$MIN_MEMORY"'))?0:1}' || exit 0

	# IMPORTANT: Let ssd do the backgrounding so we immediatly get a valid
	# pid file in the next step (ionice)
	start-stop-daemon --start --quiet --background \
		--make-pidfile --pidfile ${PIDFILE} \
		--exec @sbindir@/@PACKAGE@ -- \
		-l ${PRELOAD_LOGFILE:-/var/log/preload.log} -V ${PRELOAD_VERBOSITY:-1} \
		-n ${PRELOAD_NICE:-15} -s ${PRELOAD_STATEFILE:-""} ${PRELOAD_OPTS} -f

	IONICE=$(which ionice)
	if [ -x "$IONICE" ]; then
		IONICE_OPTS=${IONICE_OPTS:--c3}
		$IONICE ${IONICE_OPTS} -p$(<${PIDFILE})
	fi
	eend $?
}

stop() {
	ebegin "Stopping @PACKAGE@ (this may take while)"
	start-stop-daemon --stop --retry 120 --quiet --pidfile ${PIDFILE}
	eend $?
}
