# media-sound/squeezeslave/files/squeezeslave.initd-r1
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="/usr/bin/${SVCNAME}"
command_args="--daemonize ${SS_LOG} --mac ${SS_MAC} --retry ${LMS_HOST} ${SS_OPTS}"
pidfile="/var/run/${SVCNAME}.pid"
start_stop_daemon_args="--background --quiet --make-pidfile"

depend() {
	need net
	use logitechmediaserver
}
