# net-wireless/iwd/files/ead.initd
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

pidfile="/run/ead.pid"
command="/usr/libexec/ead"
command_background="yes"

depend() {
	need dbus
	keyword -shutdown
}
