# net-wireless/iwd/files/iwd.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

pidfile="/run/iwd.pid"
command="/usr/libexec/iwd"
command_background="yes"

depend() {
	need dbus
	keyword -shutdown
}
