# app-crypt/moolticute/files/moolticuted.init
#!/sbin/openrc-run
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command=/usr/bin/moolticuted
command_background="yes"
pidfile=/var/run/moolticuted.pid
command_user="nobody:usb"

depend() {
	need udev
}
