# net-vpn/i2p/files/i2p.init
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

readonly WRAPPER_CONF="/usr/share/i2p/wrapper.config"

command=wrapper
command_args="$WRAPPER_CONF wrapper.syslog.ident=i2p"
command_background=true
command_progress=true
command_user=i2p
pidfile=/var/run/i2p.pid
retry='SIGTERM/20 SIGKILL/20'
