# net-misc/bird/files/initd-v6-bird-1.3.1
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

opts="reload"

EXE=bird6
SOCK="/var/run/${EXE}.ctl"

depend() {
	need net
	use logger
}

checkconfig() {
	if [ ! -f "/etc/${EXE}.conf" ]; then
		eerror "Please create /etc/${EXE}.conf"
		return 1
	fi
	return 0
}

start() {
	checkconfig || return $?
	ebegin "Starting BIRD"
	start-stop-daemon --start --exec /usr/sbin/${EXE} -- -c "/etc/${EXE}.conf" -s "${SOCK}"
	eend $? "Failed to start BIRD"
}

stop() {
	ebegin "Stopping BIRD"
	start-stop-daemon --stop --exec /usr/sbin/${EXE}
	eend $? "Failed to stop BIRD"
}

reload() {
	service_started "${SVCNAME}" || return
	ebegin "Reloading BIRD"
	start-stop-daemon --stop --signal HUP --oknodo --exec /usr/sbin/${EXE}
	eend $? "Failed to reload BIRD"
}
