# app-admin/metalog/files/metalog.service
[Unit]
Description=System Logger Daemon
Names=syslog.service

[Service]
ExecStart=/usr/sbin/metalog
Restart=always

[Install]
Alias=syslog.service
