# app-containers/podman/files/podman-auto-update-5.0.0_rc4.cron
#!/bin/sh
set -e

# mimicking behaviour of podman-auto-update.service
if [ -x /usr/bin/podman ]; then
   /usr/bin/podman auto-update && /usr/bin/podman image prune -f
fi
