# app-containers/podman/files/podman-auto-update-5.0.0.cron
#!/bin/sh

# podman-auto-update(1)
# uncomment following to auto update containers according to their auto-update policy
# /usr/bin/podman auto-update && /usr/bin/podman image prune -f
