# sys-apps/lm_sensors/files/lm_sensors.initd
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="/usr/bin/sensors"

depend() {
	need localmount
	after modules-load
}

start() {
	ebegin "Initializing sensors"
	${command} --set >/dev/null 2>&1
	eend $?
}

stop() {
	# Nothing to stop
	return 0;
}
