# net-analyzer/rrdtool/files/rrdcached.init
#!/sbin/openrc-run
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="/usr/bin/rrdcached"
command_args="${RRCACHE_ARGS} -p /run/rrdcached.pid"
start_stop_daemon_args="--quiet"
description="RRDtool data caching daemon"

depend() {
	need localmount net
	use logger
}
