# app-admin/salt/files/api-initd-3
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

command="/usr/bin/salt-api"
command_args="${SALT_OPTS}"
command_background="1"
pidfile="/var/run/salt-api.pid"
name="SALT API daemon"

depend() {
	use net logger
}
