# app-admin/salt/files/api-initd-5
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/salt-api"
command_args="${SALT_OPTS}"
pidfile="/var/run/supervisor-salt-api.pid"
daemon_pidfile="/var/run/salt-api.pid"
command_args_background="--daemon --pid-file=\"${daemon_pidfile}\""
name="SALT API daemon"
retry="20"

supervisor=supervise-daemon
output_logger="logger -p daemon.info -t salt-api.stdout --"
error_logger="logger -p daemon.err -t salt-api.error --"

depend() {
	use net logger
}
