# app-admin/salt/files/syndic-initd-5
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/salt-syndic"
command_args="${SALT_OPTS}"
pidfile="/var/run/supervisor-salt-syndic.pid"
daemon_pidfile="/var/run/salt-syndic.pid"
command_args_background="--daemon --pid-file=\"${daemon_pidfile}\""
name="SALT syndic daemon"
retry="20"

supervisor=supervise-daemon
output_logger="logger -p daemon.info -t salt-syndic.stdout --"
error_logger="logger -p daemon.err -t salt-syndic.error --"

depend() {
	use net logger
}
