# dev-vcs/cvsd/files/cvsd.rc6
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting cvsd"
	start-stop-daemon --start --quiet --background --exec /usr/sbin/cvsd
	eend $?
}

stop() {
	ebegin "Stopping cvsd"
	start-stop-daemon --stop --quiet --pidfile /var/run/cvsd.pid
	eend $?
}
