# net-misc/utelnetd/files/utelnetd.initd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

start() {
	ebegin "Starting utelnetd"
	start-stop-daemon -b --start --quiet --exec /usr/sbin/utelnetd
	eend $?
}

stop() {
	ebegin "Stopping utelnetd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/utelnetd
	eend $?
}
