# net-misc/sks/files/sks-recon.service
[Unit]
Description=SKS key reconnaissance thread
After=network.target

[Service]
ExecStart=/usr/bin/sks recon
User=sks
WorkingDirectory=/var/lib/sks

[Install]
WantedBy=multi-user.target
