# net-p2p/syncthing/files/syncthing.initd-user
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

description="Open Source Continuous File Synchronization"
command="/usr/bin/syncthing"
supervisor=supervise-daemon

command_args="serve --no-browser --no-restart"

depend() {
	need localmount
	after net
}

start_pre() {
	checkpath -q -d -o ${command_user} ${SYNCTHING_HOMEDIR}
	checkpath -q -f -o ${command_user} ${SYNCTHING_LOGFILE}
}
