# net-p2p/syncthing/files/syncthing.initd-user
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License, v2 or later

supervisor=supervise-daemon
command="/usr/bin/syncthing"
command_args="serve --no-browser --no-restart --logflags=0"
