# app-misc/inputlircd/files/inputlircd.init.2
#!/sbin/openrc-run

depend() {
	provide lirc
}

start() {
	ebegin "Starting inputlircd"
	start-stop-daemon --start --quiet --exec /usr/sbin/inputlircd -- ${INPUTLIRCD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping inputlircd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/inputlircd
	eend $?
}
