# net-misc/oidentd/files/oidentd.service
[Unit]
Description=TCP/IP IDENT protocol server

[Service]
ExecStart=/usr/sbin/oidentd -i -S -u nobody -g nobody
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
