# net-misc/oidentd/files/oidentd.socket
[Unit]
Description=Ident (RFC 1413) socket
Conflicts=oidentd.service

[Socket]
ListenStream=113
Accept=yes

[Install]
WantedBy=sockets.target
