# app-admin/fleet/files/fleetd.service
[Unit]
Description=fleet
Wants=etcd.service

[Service]
ExecStart=/usr/bin/fleetd

[Install]
WantedBy=multi-user.target
