# app-admin/fleet/files/fleet.socket
[Socket]
ListenStream=/var/run/fleet.sock
