# sys-apps/hd-idle/files/hd-idle-init
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	after bootmisc
}

start() {
	ebegin "Starting hd-idle"
	start-stop-daemon --start --exec /usr/sbin/hd-idle -- ${HD_IDLE_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping hd-idle"
	start-stop-daemon --stop --exec /usr/sbin/hd-idle
	eend $?
}
