# sys-apps/hd-idle/files/hd-idle.service
# Copyright 1999-2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2
[Unit]
Description=hd-idle - spinning down HDDs after a period of idle time
Documentation=man:hd-idle(1)

[Service]
Type=forking
ExecStart=/usr/sbin/hd-idle -i 180 -l /var/log/hd-idle.log

[Install]
WantedBy=multi-user.target
