# net-misc/dhcpcd/files/dhcpcd.initd
#!/sbin/openrc-run
# Copyright 2007-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command=/sbin/dhcpcd
pidfile=/var/run/dhcpcd.pid
command_args=-q
name="DHCP Client Daemon"

depend()
{
	provide net
	need localmount
	use logger network
	after bootmisc modules
	before dns
}
