# net-dns/pdns-recursor/files/pdns-recursor-r3
#!/sbin/openrc-run
# Copyright 1999-2024 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

PDNS_CONFIGDIR="${PDNS_CONFIGDIR:-/etc/powerdns}"
[ "${RC_SVCNAME}" != "pdns-recursor" ] && PDNS_INSTANCE="${RC_SVCNAME#pdns-recursor.}" || PDNS_INSTANCE=""

name="PowerDNS Recursor"
description="Recursive name server"

extra_started_commands="ping"
description_ping="Ping the PowerDNS instance"

command="/usr/sbin/pdns_recursor"
command_args="--config-dir=${PDNS_CONFIGDIR}${PDNS_INSTANCE:+ --config-name=${PDNS_INSTANCE}} --write-pid=yes"
command_args_foreground="--daemon=no"
command_args_background="--daemon=yes"
yesno ${rc_verbose} || command_args_background="${command_args_background} >/dev/null 2>&1"

control_command="/usr/bin/rec_control"
control_command_args="${command_args}"

depend() {
	need net
}

ping() {
	ebegin "Pinging ${name}"
	${control_command} ${control_command_args} ping >/dev/null 2>&1
	eend $?
}
