# net-analyzer/zabbix/files/zabbix-agentd.service
[Unit]
Description=Zabbix agent daemon
After=network.target

[Service]
Type=forking
User=zabbix
Group=zabbix
ExecStart=/usr/sbin/zabbix_agentd -c /etc/zabbix/zabbix_agentd.conf
PIDFile=/run/zabbix/zabbix_agentd.pid

[Install]
WantedBy=multi-user.target

