# net-analyzer/zabbix/files/zabbix-proxy.service
[Unit]
Description=Zabbix proxy
After=network.target

[Service]
Type=forking
User=zabbix
Group=zabbix
ExecStart=/usr/sbin/zabbix_proxy -c /etc/zabbix/zabbix_proxy.conf
PIDFile=/run/zabbix/zabbix_proxy.pid

[Install]
WantedBy=multi-user.target

