# media-radio/ax25-tools/files/ax25d.rc
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/media-radio/ax25-tools/files/ax25d.rc,v 1.5 2011/01/08 17:12:50 tomjbe Exp $

depend() {
	need net
}

start() {
	ebegin "Starting ax25d"
	start-stop-daemon --start --quiet --background --exec /usr/sbin/ax25d
	eend $?
}

stop() {
	ebegin "Stopping ax25d"
	start-stop-daemon --stop --quiet --exec /usr/sbin/ax25d
	eend $?
}
