# media-radio/ax25-tools/files/rxecho.rc
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	need net
}

start() {
	ebegin "Starting rxecho"
	start-stop-daemon --start --quiet --background --exec /usr/sbin/rxecho
	eend $?
}

stop() {
	ebegin "Stopping rxecho"
	start-stop-daemon --stop --quiet --exec /usr/sbin/rxecho
	eend $?
}
