# x11-terms/rxvt-unicode/files/urxvtd.socket
[Unit]
Description=rxvt-unicode terminal daemon socket

[Socket]
ListenStream=%h/.urxvt/urxvtd-%H

[Install]
WantedBy=sockets.target
