# x11-terms/rxvt-unicode/files/urxvtd.service
[Unit]
Description=rxvt-unicode terminal daemon
Documentation=man:urxvtd

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=default.target
