# sys-cluster/sanlock/files/sanlock.initd
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

pidfile="/var/run/snlock.pid"
command="/usr/sbin/snlock"
command_args="daemon ${SANLOCKOPTS:-"-U sanlock -G sanlock"}"

