# sys-apps/rng-tools/files/rngd-initd-4.1
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-apps/rng-tools/files/rngd-initd-4.1,v 1.7 2013/05/22 05:16:17 vapier Exp $

depend() {
	need localmount
	after urandom
	provide entropy
}

# Do NOT add /dev/tpm to this.
DEFAULT_DEVICE="/dev/hw_random* /dev/hwrandom* /dev/i810_rng /dev/hwrng*"

find_device() {
	local d
	# The echo is to cause globbing
	for d in $(echo ${DEFAULT_DEVICE}) ; do
		[ -e "${d}" ] && echo "${d}"
	done
}

RNG_DEVICE="${DEVICE:-$(find_device)}"

command=/usr/sbin/rngd
pidfile="/var/run/${SVCNAME}.pid"
command_args="--pid-file ${pidfile} --background --random-step ${STEP:-64} ${NO_TPM:+--no-tpm=1} ${NO_DRNG:+--no-drng=1} --fill-watermark ${WATERMARK} ${RNG_DEVICE:+--rng-device ${RNG_DEVICE}}"
start_stop_daemon_args="--wait 1000"
retry="SIGKILL/5000"
