# sys-apps/rng-tools/files/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
ExecStart=/usr/sbin/rngd -f

[Install]
WantedBy=multi-user.target
