# net-dns/nsd/files/nsd3-patch.cron
#!/bin/sh
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/nsd/files/nsd3-patch.cron,v 1.1 2013/04/03 19:10:39 wschlich Exp $

# Cron job: merge zone transfer changes back to zone files
# Needed for nsd3 only

/etc/init.d/nsd --quiet status && /etc/init.d/nsd patch >/dev/null
