# dev-python/mwlib/files/nslave.initd-r1
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="nslave pulls new jobs from a mw-qserve instance"

command="/usr/bin/nslave"
command_args="--cachedir=${NSLAVE_CACHEDIR} \
	--serve-files-address=${NSLAVE_ADDRESS} \
	--serve-files-port=${NSLAVE_PORT} \
	--url=${NSLAVE_URL} \
	--numprocs=${NSLAVE_NUMPROCS}"

pidfile="/var/run/${SVCNAME}.pid"
command_background="yes"

start_stop_daemon_args="--user ${NSLAVE_USER:-mwlib} \
	--group ${NSLAVE_GROUP:-mwlib} \
	${NSLAVE_LOGFILE:+--stderr ${NSLAVE_LOGFILE}} \
	${NSLAVE_LOGFILE:+--stdout ${NSLAVE_LOGFILE}} \
	--umask ${NSLAVE_UMASK}"

depend() {
	need net localmount
	use mw-qserve
}
