# dev-python/mwlib/files/postman.initd-r1
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="postman uploads zip collections to pediapress in case someone likes to order printed books"

command="/usr/bin/postman"
command_args="--cachedir=${POSTMAN_CACHEDIR}"

pidfile="/var/run/${SVCNAME}.pid"
command_background="yes"

start_stop_daemon_args="--user ${POSTMAN_USER:-mwlib} \
	--group ${POSTMAN_GROUP:-mwlib} \
	${POSTMAN_LOGFILE:+--stderr ${POSTMAN_LOGFILE}} \
	${POSTMAN_LOGFILE:+--stdout ${POSTMAN_LOGFILE}}"

depend() {
	need net localmount
}
