# dev-db/mongodb/files/mongos.initd-r2
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-db/mongodb/files/mongos.initd-r2,v 1.2 2014/05/26 08:03:45 ultrabug Exp $

config_file=${config_file:-/etc/${SVCNAME}.conf}
run_dir=${run_dir:-/run/mongodb}

command="/usr/bin/mongos"
command_args="--config ${config_file}"
command_background="true"
pidfile=${run_dir}/${SVCNAME}.pid
user=${user:-mongodb}
group=${group:-mongodb}
start_stop_daemon_args="--user ${user} --group ${group}"

depend() {
	use net
}

start_pre() {
    checkpath -d -m 0750 -o "${user}":"${group}" "${run_dir}"
    if [ ! -f ${config_file} ]; then
    	eerror "Missing configuration file ${config_file}"
    	return 1
    fi
}
