# dev-db/mongodb/files/mongos.initd-r3
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

: ${config_file:="/etc/${RC_SVCNAME}.conf"}
: ${user:=mongodb}
: ${group:=mongodb}

command="/usr/bin/mongos"
command_args="--config ${config_file}"
command_background="true"
pidfile="/run/${RC_SVCNAME}.pid"
command_user="${user}:${group}"
required_files="${config_file}"

depend() {
	use net
}
