# net-dialup/dwun/files/dwun
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dialup/dwun/files/dwun,v 1.4 2014/12/04 08:27:53 pinkbyte Exp $

depend() {
	need net
}

DAEMON=/usr/sbin/dwun
PIDFILE="/run/${SVCNAME}.pid"

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --pidfile ${PIDFILE} --startas ${DAEMON}
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --signal 1 --pidfile ${PIDFILE}
	eend $?
}
