# app-crypt/dehydrated/files/dehydrated.timer
[Unit]
Description=Run dehydrated ACME client

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3h

[Install]
WantedBy=timers.target
