# net-misc/asterisk/files/asterisk.service
[Unit]
Description=Asterisk PBX and telephony daemon.
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/lib/asterisk
User=asterisk
Group=asterisk
ExecStart=/usr/sbin/asterisk -f -C /etc/asterisk/asterisk.conf
ExecStop=/usr/sbin/asterisk -rx 'core stop now'
ExecReload=/usr/sbin/asterisk -rx 'core reload'

LimitNOFILE=4096
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
