# net-mail/qmail-notify/files/qmail-notify.cron
#!/bin/sh
#
# uncomment this to make it work
#

#exec /usr/sbin/qmail-notify -r
