# net-ftp/vsftpd/files/vsftpd.socket
[Unit]
Conflicts=vsftpd.service

[Socket]
ListenStream=21
Accept=yes

[Install]
WantedBy=sockets.target
